# apache2-minimal/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio

# Enable strict mode
set -e
# shellcheck disable=SC1091



echo "Starting Apache2..."

ssl=$(bashio::config 'ssl')
# Get local IP
# We can use the Supervisor API to get the Host IP
host_ip=""
if [ -n "$SUPERVISOR_TOKEN" ]; then
    network_info=$(curl -sSL -H "Authorization: Bearer $SUPERVISOR_TOKEN" http://supervisor/network/info 2>/dev/null)
    if [ -n "$network_info" ] && echo "$network_info" | jq -e . >/dev/null 2>&1; then
        host_ip=$(echo "$network_info" | jq -r '.data.interfaces[]? | select(.primary == true) | .ipv4.address[]?' 2>/dev/null | head -n 1 | cut -d'/' -f1)
    fi
fi

# Fallback if Host IP is empty
if [ -z "$host_ip" ]; then
    host_ip="homeassistant.local"
fi

# Determine protocol and port
protocol="http"
port="80"
if [ "$ssl" = "true" ]; then
    protocol="https"
    port="443"
fi

# Get External Port from Supervisor
external_port=""
if [ -n "$SUPERVISOR_TOKEN" ]; then
    addon_info=$(curl -sSL -H "Authorization: Bearer $SUPERVISOR_TOKEN" http://supervisor/addons/self/info 2>/dev/null)
    if [ -n "$addon_info" ] && echo "$addon_info" | jq -e . >/dev/null 2>&1; then
        external_port=$(echo "$addon_info" | jq -r ".data.network.\"${port}/tcp\"?" 2>/dev/null)
    fi
fi

bashio::log.info "---------------------------------------------------"
bashio::log.info "Apache2 is starting..."
bashio::log.info ""
if [ "$external_port" != "null" ] && [ -n "$external_port" ]; then
    bashio::log.info "You can access the website in your local network at:"
    bashio::log.info "  ${protocol}://${host_ip}:${external_port}"
else
    bashio::log.warning "Could not determine external port mapping."
    bashio::log.warning "Please check your port configuration in Home Assistant."
    bashio::log.info "Internal container connection: ${protocol}://$(hostname -i):${port}"
fi
bashio::log.info ""
bashio::log.info "---------------------------------------------------"
exec /usr/sbin/httpd -D FOREGROUND