# wordpress/rootfs/etc/s6-overlay/s6-rc.d/init-wordpress/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Initialize Wordpress..."

# SSL Configuration
SSL_CERT="/ssl/$(bashio::config 'certfile')"
SSL_CERT_KEY="/ssl/$(bashio::config 'keyfile')"
export SSL_CERT SSL_CERT_KEY

# Configure database connection
if bashio::services.available "mysql"; then
    bashio::log.info "Configuring Wordpress to use the MySQL service..."
    DB_HOST=$(bashio::services "mysql" "host")
    DB_PORT=$(bashio::services "mysql" "port")
    DB_USER=$(bashio::services "mysql" "username")
    DB_PASS=$(bashio::services "mysql" "password")
    DB_NAME="wordpress"

    # Combine host and port for Wordpress config
    if bashio::var.has_value "${DB_PORT}"; then
        DB_HOST_CONFIG="${DB_HOST}:${DB_PORT}"
    else
        DB_HOST_CONFIG="${DB_HOST}"
        DB_PORT="3306"
    fi
else
    bashio::exit.nok "The 'mysql' service is not available! Please install and start the MariaDB add-on."
fi

# Removed insecure password check as we now auto-generate it


    # Ensure database exists and is accessible
    bashio::log.info "Checking database connection..."
    bashio::log.info "Database Host: ${DB_HOST}"
    bashio::log.info "Database Port: ${DB_PORT}"

    # Create temporary config file for secure authentication
    temp_config=$(mktemp)
    # Ensure cleanup on any exit
    trap 'rm -f "$temp_config"' EXIT
    install -m 600 /dev/null "$temp_config"
    # Use printf for safe credential handling (avoids shell interpolation issues)
    printf '[client]\npassword=%s\n' "$DB_PASS" > "$temp_config"

    # Wait for loop
    for i in {1..30}; do
        if mariadb-admin --defaults-extra-file="$temp_config" ping --protocol=tcp --skip-ssl -h "$DB_HOST" -P "$DB_PORT" -u "$DB_USER" --silent; then
            break
        fi
        bashio::log.info "Waiting for database... $i/30"
        sleep 2
    done

    # Verify connection after loop (removing --silent to see error on last failure)
    if ! mariadb-admin --defaults-extra-file="$temp_config" ping --protocol=tcp --skip-ssl -h "$DB_HOST" -P "$DB_PORT" -u "$DB_USER"; then
        bashio::log.error "Could not connect to database after retries!"
        rm -f "$temp_config"
        exit 1
    fi

    # Create DB
    bashio::log.info "Ensuring database '$DB_NAME' exists..."
    if ! mariadb --defaults-extra-file="$temp_config" --protocol=tcp --skip-ssl -h "$DB_HOST" -P "$DB_PORT" -u "$DB_USER" -e "USE \`$DB_NAME\`" 2>/dev/null; then
         bashio::log.info "Database '$DB_NAME' not found. Creating..."
         if ! mariadb --defaults-extra-file="$temp_config" --protocol=tcp --skip-ssl -h "$DB_HOST" -P "$DB_PORT" -u "$DB_USER" -e "CREATE DATABASE IF NOT EXISTS \`$DB_NAME\` CHARACTER SET utf8mb4 COLLATE utf8mb4_unicode_ci"; then
            bashio::log.error "Failed to create database '$DB_NAME'!"
            rm -f "$temp_config"
            exit 1
         fi
         bashio::log.info "Database created."
    fi

    # Cleanup
    rm -f "$temp_config"

# Check if Wordpress is installed
if [ ! -d /var/www/wordpress/wp-admin ]; then
    bashio::log.info "Installing Wordpress..."
    if [ ! -d /usr/src/wordpress ]; then
        bashio::exit.nok "Source directory /usr/src/wordpress not found!"
    fi
    cp -r /usr/src/wordpress/. /var/www/wordpress || bashio::exit.nok "Failed to copy Wordpress files!"
    bashio::log.info "Wordpress installed!"
fi

# Configuration
bashio::log.info "Updating Wordpress configuration..."

# Get WordPress URL for HTTP_HOST
WORDPRESS_URL_CONFIG=$(bashio::config 'wordpress_url')
if [ -z "$WORDPRESS_URL_CONFIG" ] || [ "$WORDPRESS_URL_CONFIG" == "null" ]; then
    WORDPRESS_URL_CONFIG="http://wordpress.local"
fi

# Extract host from URL for HTTP_HOST
HTTP_HOST_VAL="${WORDPRESS_URL_CONFIG#http://}"
HTTP_HOST_VAL="${HTTP_HOST_VAL#https://}"
HTTP_HOST_VAL="${HTTP_HOST_VAL%%/*}"

# Check if wp-config.php exists, if not create it
if [ ! -f /var/www/wordpress/wp-config.php ]; then
    bashio::log.info "Creating wp-config.php..."
    if ! wp config create \
        --allow-root \
        --path=/var/www/wordpress \
        --dbname="${DB_NAME}" \
        --dbuser="${DB_USER}" \
        --dbpass="${DB_PASS}" \
        --dbhost="${DB_HOST_CONFIG}" \
        --skip-check \
        --force; then
            bashio::exit.nok "Failed to create wp-config.php!"
    fi

    # Add HTTP_HOST fix to wp-config.php to avoid warnings during CLI operations
    if ! grep -q "Set HTTP_HOST for CLI" /var/www/wordpress/wp-config.php 2>/dev/null; then
        # Find the line with "That's all, stop editing!" and insert before it
        sed -i "/That's all, stop editing/i\\
/* Set HTTP_HOST for CLI operations to avoid warnings */\\
if (defined('WP_CLI') && WP_CLI && !isset(\$_SERVER['HTTP_HOST'])) {\\
    \$_SERVER['HTTP_HOST'] = '${HTTP_HOST_VAL}';\\
}\\
" /var/www/wordpress/wp-config.php
    fi
else
    bashio::log.info "Updating existing wp-config.php with current database details..."
    # Update existing config with current values (in case IP or password changed)
    wp config set DB_NAME "${DB_NAME}" --allow-root --path=/var/www/wordpress
    wp config set DB_USER "${DB_USER}" --allow-root --path=/var/www/wordpress
    wp config set DB_PASSWORD "${DB_PASS}" --allow-root --path=/var/www/wordpress
    wp config set DB_HOST "${DB_HOST_CONFIG}" --allow-root --path=/var/www/wordpress

    # Ensure HTTP_HOST fix is in wp-config.php
    if ! grep -q "Set HTTP_HOST for CLI" /var/www/wordpress/wp-config.php 2>/dev/null; then
        sed -i "/That's all, stop editing/i\\
/* Set HTTP_HOST for CLI operations to avoid warnings */\\
if (defined('WP_CLI') && WP_CLI && !isset(\$_SERVER['HTTP_HOST'])) {\\
    \$_SERVER['HTTP_HOST'] = '${HTTP_HOST_VAL}';\\
}\\
" /var/www/wordpress/wp-config.php
    fi
fi

# Salts are automatically generated by 'wp config create'.
# If config existed, salts are preserved.
# We don't need to manually fetch them anymore.

# Permissions
bashio::log.info "Fixing permissions..."
chown -R nginx:nginx /var/www/wordpress || bashio::exit.nok "Failed to set ownership for /var/www/wordpress"
find /var/www/wordpress -type d -exec chmod 755 {} \; || bashio::exit.nok "Failed to set directory permissions"
find /var/www/wordpress -type f -exec chmod 644 {} \; || bashio::exit.nok "Failed to set file permissions"
chmod 640 /var/www/wordpress/wp-config.php || bashio::exit.nok "Failed to set permissions for wp-config.php"

bashio::log.info "Wordpress initialized!"

# Check if Wordpress is already installed
if ! wp --allow-root --path=/var/www/wordpress core is-installed; then
    bashio::log.info "Wordpress not installed. Installing..."

    # Generate secure password
    ADMIN_PASS=$(openssl rand -base64 24)
    ADMIN_USER=$(bashio::config 'wordpress_admin_user')
    ADMIN_EMAIL=$(bashio::config 'wordpress_admin_email')
    TITLE=$(bashio::config 'wordpress_title')
    WORDPRESS_URL=$(bashio::config 'wordpress_url')

    if [ -z "$WORDPRESS_URL" ] || [ "$WORDPRESS_URL" == "null" ]; then
        bashio::exit.nok "The configuration 'wordpress_url' is missing! Please set your external URL (e.g. https://example.com) in the add-on configuration."
    fi

    # Validate required configuration
    MISSING_CONFIGS=""
    if [ -z "$ADMIN_USER" ] || [ "$ADMIN_USER" == "null" ]; then
        MISSING_CONFIGS="${MISSING_CONFIGS}wordpress_admin_user "
    fi
    if [ -z "$ADMIN_EMAIL" ] || [ "$ADMIN_EMAIL" == "null" ]; then
        MISSING_CONFIGS="${MISSING_CONFIGS}wordpress_admin_email "
    fi
    if [ -z "$TITLE" ] || [ "$TITLE" == "null" ]; then
        MISSING_CONFIGS="${MISSING_CONFIGS}wordpress_title "
    fi
    if [ -n "$MISSING_CONFIGS" ]; then
        bashio::log.error "Missing required configuration: ${MISSING_CONFIGS}"
        bashio::log.error "Please configure these options in the add-on settings."
        exit 1
    fi

    # Set HTTP_HOST environment variable to avoid warnings during CLI installation
    # Extract host from URL
    HTTP_HOST_VAR="${WORDPRESS_URL#http://}"
    HTTP_HOST_VAR="${HTTP_HOST_VAR#https://}"
    HTTP_HOST_VAR="${HTTP_HOST_VAR%%/*}"
    export HTTP_HOST="${HTTP_HOST_VAR}"
    export REQUEST_URI="/"

    # Install Wordpress
    # Note: The HTTP_HOST warning is harmless and occurs during CLI installation
    # It doesn't affect the installation process
    if wp core install \
        --allow-root \
        --path=/var/www/wordpress \
        --url="${WORDPRESS_URL}" \
        --title="${TITLE}" \
        --admin_user="${ADMIN_USER}" \
        --admin_password="${ADMIN_PASS}" \
        --admin_email="${ADMIN_EMAIL}"; then

        bashio::log.info "Wordpress installed successfully!"
        bashio::log.warning "------------------------------------------------"
        # Security Note: We are intentionally logging the password here for user convenience
        # during the initial setup, as requested by the user.
        bashio::log.warning "Wordpress Admin Password: ${ADMIN_PASS}"
        bashio::log.warning "PLEASE SAVE THIS PASSWORD NOW!"
        bashio::log.warning "This password will NOT be shown again."
        bashio::log.warning "------------------------------------------------"
    else
        bashio::log.error "Failed to install Wordpress!"
    fi
else
    bashio::log.info "Wordpress database already installed."
fi

# Configure nginx based on SSL settings
bashio::log.info "Configuring nginx..."

# Remove old ingress.conf to avoid conflicts
if [ -f /etc/nginx/http.d/ingress.conf ]; then
	rm -f /etc/nginx/http.d/ingress.conf
	bashio::log.info "Removed old ingress.conf"
fi

# Create .well-known directory for Let's Encrypt if using HTTPS
mkdir -p /var/www/wordpress/.well-known

# Checks if SSL certificate and key exists, otherwise default to http traffic
if bashio::config.true 'ssl'; then
	if [ ! -f "${SSL_CERT}" ] || [ ! -f "${SSL_CERT_KEY}" ]; then
		bashio::log.error "SSL disabled: Certificate or key file not found!"
		bashio::log.error "Make sure ${SSL_CERT} and ${SSL_CERT_KEY} exist."
		bashio::log.error "Reverting to HTTP."
		cat /etc/nginx/templates/http.conf >/etc/nginx/http.d/default.conf
	else
		bashio::log.info "SSL has been enabled. Setting nginx settings for SSL usage with ${SSL_CERT} and ${SSL_CERT_KEY}."
		# Use envsubst to replace SSL_CERT and SSL_CERT_KEY variables
		export SSL_CERT SSL_CERT_KEY
		envsubst '${SSL_CERT} ${SSL_CERT_KEY}' </etc/nginx/templates/https.conf >/etc/nginx/http.d/default.conf
	fi
else
	bashio::log.info "SSL is disabled, using HTTP."
	cat /etc/nginx/templates/http.conf >/etc/nginx/http.d/default.conf
fi

# Verify nginx configuration
if nginx -t 2>&1; then
	bashio::log.info "Nginx configuration is valid."
	bashio::log.info "Nginx configuration file created at /etc/nginx/http.d/default.conf"
else
	bashio::log.error "Nginx configuration test failed!"
	bashio::log.error "Configuration file content:"
	cat /etc/nginx/http.d/default.conf || true
	exit 1
fi

# Log access information
if bashio::config.true 'ssl'; then
	bashio::log.info "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
	bashio::log.info "WordPress is accessible via:"
	bashio::log.info "  HTTPS: https://<your-home-assistant-ip>:8449"
	bashio::log.info "  HTTP:  http://<your-home-assistant-ip>:8099 (redirects to HTTPS)"
	bashio::log.info "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
else
	bashio::log.info "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
	bashio::log.info "WordPress is accessible via:"
	bashio::log.info "  HTTP:  http://<your-home-assistant-ip>:8099"
	bashio::log.info "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
fi

# Warning if wordpress_url doesn't match SSL setting
if bashio::config.true 'ssl' && [[ "${WORDPRESS_URL}" == http://* ]]; then
	bashio::log.warning "SSL is enabled, but wordpress_url starts with http://. Consider changing it to https://."
elif ! bashio::config.true 'ssl' && [[ "${WORDPRESS_URL}" == https://* ]]; then
	bashio::log.warning "SSL is disabled, but wordpress_url starts with https://. Consider changing it to http://."
fi
