# wordpress/rootfs/etc/s6-overlay/s6-rc.d/init-wordpress/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Initialize Wordpress..."

# Set defaults for optional values
if ! bashio::config.has_value 'database_host'; then
    if bashio::services.available "mysql"; then
        bashio::log.info "Using MySQL service connection details..."
        DB_HOST=$(bashio::services "mysql" "host")
        if bashio::var.has_value "$(bashio::services.port "mysql")"; then
            DB_HOST="${DB_HOST}:$(bashio::services.port "mysql")"
        fi
    else
        bashio::log.info "No database_host configured, using Home Assistant MariaDB addon (core-mariadb)"
        DB_HOST="core-mariadb"
    fi
else
    DB_HOST=$(bashio::config 'database_host')
fi

if ! bashio::config.has_value 'database_password'; then
    if bashio::services.available "mysql"; then
        DB_PASS=$(bashio::services.credentials "mysql" "password")
    else
        bashio::exit.nok "A 'database_password' is required!"
    fi
else
    DB_PASS=$(bashio::config 'database_password')
fi

if ! bashio::config.has_value 'database_name'; then
    bashio::log.info "No database_name configured, using default 'wordpress'"
    DB_NAME="wordpress"
else
    DB_NAME=$(bashio::config 'database_name')
fi

if ! bashio::config.has_value 'database_user'; then
    if bashio::services.available "mysql"; then
        DB_USER=$(bashio::services.credentials "mysql" "username")
    else
        bashio::log.info "No database_user configured, using default 'wordpress'"
        DB_USER="wordpress"
    fi
else
    DB_USER=$(bashio::config 'database_user')
fi

if bashio::config.equals 'wordpress_admin_password' 'changeme'; then
    bashio::exit.nok "You are using the insecure default admin password 'changeme'. Please change it!"
fi

# Check if Wordpress is installed
if [ ! -d /var/www/wordpress/wp-admin ]; then
    bashio::log.info "Installing Wordpress..."
    if [ ! -d /usr/src/wordpress ]; then
        bashio::exit.nok "Source directory /usr/src/wordpress not found!"
    fi
    cp -r /usr/src/wordpress/. /var/www/wordpress || bashio::exit.nok "Failed to copy Wordpress files!"
    bashio::log.info "Wordpress installed!"
fi

# Configuration
if [ ! -f /var/www/wordpress/wp-config.php ]; then
    bashio::log.info "Creating configuration..."
    mv /var/www/wordpress/wp-config-sample.php /var/www/wordpress/wp-config.php || bashio::exit.nok "Failed to create wp-config.php!"

    # Database - use variables set during validation
    # Escape special characters for sed
    # We escape backslashes first, then the delimiter (|) and ampersands
    db_host_escaped=$(echo "${DB_HOST}" | sed 's/\\/\\\\/g' | sed 's/[|&]/\\&/g')
    db_name_escaped=$(echo "${DB_NAME}" | sed 's/\\/\\\\/g' | sed 's/[|&]/\\&/g')
    db_user_escaped=$(echo "${DB_USER}" | sed 's/\\/\\\\/g' | sed 's/[|&]/\\&/g')
    db_pass_escaped=$(echo "${DB_PASS}" | sed 's/\\/\\\\/g' | sed 's/[|&]/\\&/g')

    sed -i "s|database_name_here|${db_name_escaped}|" /var/www/wordpress/wp-config.php
    sed -i "s|username_here|${db_user_escaped}|" /var/www/wordpress/wp-config.php
    sed -i "s|password_here|${db_pass_escaped}|" /var/www/wordpress/wp-config.php
    sed -i "s|localhost|${db_host_escaped}|" /var/www/wordpress/wp-config.php

    # Salt
    bashio::log.info "Generating salts..."
    if ! curl -s https://api.wordpress.org/secret-key/1.1/salt/ > /tmp/wp-salts; then
         bashio::log.warning "Failed to fetch unique keys and salts from WordPress.org! Using defaults (potentially insecure)."
         rm -f /tmp/wp-salts
    else
        # Validate salt content
        if [ ! -s /tmp/wp-salts ] || ! grep -q "define('AUTH_KEY'" /tmp/wp-salts; then
            bashio::log.warning "Downloaded salt file is empty or invalid! Using defaults (potentially insecure)."
            rm -f /tmp/wp-salts
        else
            # Request successful and valid: Replace placeholder salts
            # Delete the default salt definitions and append new ones.
            sed -i "/define('AUTH_KEY'/,/define('NONCE_SALT'/d" /var/www/wordpress/wp-config.php
            cat /tmp/wp-salts >> /var/www/wordpress/wp-config.php
            rm -f /tmp/wp-salts
            bashio::log.info "Authentication keys and salts successfully updated."
        fi
    fi
fi

# Permissions
bashio::log.info "Fixing permissions..."
chown -R nginx:nginx /var/www/wordpress || bashio::exit.nok "Failed to set ownership for /var/www/wordpress"
find /var/www/wordpress -type d -exec chmod 755 {} \; || bashio::exit.nok "Failed to set directory permissions"
find /var/www/wordpress -type f -exec chmod 644 {} \; || bashio::exit.nok "Failed to set file permissions"
chmod 640 /var/www/wordpress/wp-config.php || bashio::exit.nok "Failed to set permissions for wp-config.php"

bashio::log.info "Wordpress initialized!"
