# antigravity-server/rootfs/etc/s6-overlay/s6-rc.d/novnc/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Antigravity-Server
# Runs the NoVNC websocket proxy
# ==============================================================================

bashio::log.info "Starting NoVNC on port 6080..."

# Wait a moment for VNC to be ready
# Wait for VNC to be ready (max 10s)
# Wait for VNC to be ready (max 10s)
for i in {1..20}; do
    if python3 -c "import socket; s = socket.socket(); s.connect(('127.0.0.1', 5901))" > /dev/null 2>&1; then
        bashio::log.info "VNC server is ready."
        break
    fi
    if [ "$i" -eq 20 ]; then
        bashio::log.error "Timeout waiting for VNC server to become ready!"
        exit 1
    fi
    sleep 0.5
done

# Start NoVNC websocket proxy
exec /usr/share/novnc/utils/novnc_proxy \
    --listen 6080 \
    --vnc localhost:5901
