# apache2/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio

# Enable strict mode
set -e
# shellcheck disable=SC1091
source /usr/lib/bashio/banner.sh
bashio::addon.print_banner

echo "Starting Apache2..."
exec /usr/sbin/httpd -D FOREGROUND