# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/svc-postgres/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting PostgreSQL..."
exec su-exec postgres postgres -D /var/lib/postgresql/data
