# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/init-tt-rss/up
#!/command/execlineb
foreground { /etc/s6-overlay/s6-rc.d/init-tt-rss/run }
