# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/init-secrets/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Checking for secrets..."

SECRETS_FILE="/data/secrets.json"

# Check if secrets file exists
if ! bashio::fs.file_exists "${SECRETS_FILE}"; then
    bashio::log.info "Generating new secrets..."

    # Generate secrets
    SECRET_KEY=$(openssl rand -hex 32)
    POSTGRES_PASSWORD=$(openssl rand -hex 16)

    # Save to file with restrictive permissions
    umask 077
    bashio::var.json \
        SECRET_KEY "${SECRET_KEY}" \
        POSTGRES_PASSWORD "${POSTGRES_PASSWORD}" \
        > "${SECRETS_FILE}"
    chmod 600 "${SECRETS_FILE}"

    bashio::log.info "Secrets generated and saved."
else
    bashio::log.info "Loading existing secrets."
fi

# Export to S6 environment for other services (use -r for raw unquoted output)
SECRET_KEY=$(jq -r ".SECRET_KEY" "${SECRETS_FILE}")
POSTGRES_PASSWORD=$(jq -r ".POSTGRES_PASSWORD" "${SECRETS_FILE}")

printf "%s" "${SECRET_KEY}" > /var/run/s6/container_environment/SECRET_KEY
printf "%s" "${POSTGRES_PASSWORD}" > /var/run/s6/container_environment/POSTGRES_PASSWORD
