# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/tt-rss-updater/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting TT-RSS Updater..."
exec sudo -u nginx php83 /var/www/tt-rss/update_daemon2.php
