# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/tt-rss-updater/run
#!/usr/bin/env bashio
# shellcheck shell=bash

bashio::log.info "Starting TT-RSS Updater..."
exec s6-setuidgid nginx php84 /var/www/tt-rss/update_daemon2.php
