# netboot-xyz/root/etc/services.d/webapp/run
#!/usr/bin/with-contenv bashio
# Enable strict mode
set -e
# shellcheck disable=SC1091
source /usr/lib/bashio/banner.sh
bashio::addon.print_banner
# Get Addon Version
addon_version=$(bashio::addon.version)

# Run App
cd /app
export WEB_APP_PORT=$(bashio::addon.port 3000)
if bashio::var.has_value "$(bashio::addon.ingress_entry)"; then
    export SUBFOLDER=$(bashio::addon.ingress_entry)
fi
exec \
/usr/bin/node app.js