# netboot-xyz/root/etc/services.d/webapp/run
#!/usr/bin/with-contenv bash

# Run App
cd /app
exec \
/usr/bin/node app.js
