# netboot-xyz/root/etc/services.d/webapp/run
#!/usr/bin/with-contenv bashio
# Enable strict mode
set -e
# shellcheck disable=SC1091
source /usr/lib/bashio/banner.sh
bashio::addon.print_banner
# Get Addon Version
addon_version=$(bashio::addon.version)

# Run App
cd /app
exec \
/usr/bin/node app.js