# wordpress/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Starting Nginx..."
exec nginx -g "daemon off;"
