# wordpress/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Wait a moment for init to complete
sleep 2

# Verify nginx configuration file exists
if [ ! -f /etc/nginx/http.d/default.conf ]; then
	bashio::log.error "Nginx configuration file not found at /etc/nginx/http.d/default.conf!"
	bashio::log.error "This should have been created during initialization."
	exit 1
fi

# Test nginx configuration
if ! nginx -t 2>&1; then
	bashio::log.error "Nginx configuration test failed!"
	bashio::log.error "Configuration file content:"
	cat /etc/nginx/http.d/default.conf || true
	exit 1
fi

# Check if port 80 is accessible
if command -v ss >/dev/null 2>&1; then
	if ss -tuln | grep -q ":80 "; then
		bashio::log.warning "Port 80 is already in use!"
		ss -tuln | grep ":80 " || true
	fi
fi

bashio::log.info "Starting Nginx..."
bashio::log.info "Nginx will listen on port 80 (mapped to external port 8099)"

# Verify WordPress directory exists and is readable
if [ ! -d /var/www/wordpress ]; then
	bashio::log.error "WordPress directory /var/www/wordpress does not exist!"
	exit 1
fi

if [ ! -f /var/www/wordpress/index.php ]; then
	bashio::log.warning "WordPress index.php not found in /var/www/wordpress"
	bashio::log.warning "This might be normal if WordPress is not yet installed"
fi

# Check if PHP-FPM is running (wait up to 15 seconds)
bashio::log.info "Waiting for PHP-FPM to be ready..."
PHP_FPM_READY=false
for i in {1..15}; do
	# Check multiple ways: process name, port 9000, socket file
	if pgrep -f "php-fpm84.*master" >/dev/null 2>&1 || \
	   pgrep -f "php-fpm.*master" >/dev/null 2>&1 || \
	   (command -v ss >/dev/null 2>&1 && ss -tuln | grep -q ":9000 ") || \
	   (command -v netstat >/dev/null 2>&1 && netstat -tuln | grep -q ":9000 ") || \
	   [ -S /run/php-fpm84.sock ] 2>/dev/null || \
	   [ -S /var/run/php-fpm84.sock ] 2>/dev/null; then
		bashio::log.info "✓ PHP-FPM is running (checked: $i/15)"
		PHP_FPM_READY=true
		break
	fi
	if [ $i -lt 15 ]; then
		bashio::log.debug "PHP-FPM not ready yet, waiting... ($i/15)"
		sleep 1
	fi
done

if [ "$PHP_FPM_READY" = "false" ]; then
	bashio::log.warning "PHP-FPM not detected after 15 seconds, but continuing with Nginx startup..."
	bashio::log.warning "Nginx will start, but PHP requests may fail until PHP-FPM is ready."
fi

# Show nginx configuration summary
bashio::log.info "Nginx configuration summary:"
bashio::log.info "  - Listening on: 0.0.0.0:80 (all interfaces)"
bashio::log.info "  - External port: 8099"
bashio::log.info "  - Document root: /var/www/wordpress"
bashio::log.info "  - PHP-FPM: 127.0.0.1:9000"
bashio::log.info ""

# Verify the default.conf file exists
if [ -f /etc/nginx/http.d/default.conf ]; then
	bashio::log.info "Nginx configuration file exists and is ready"
else
	bashio::log.error "Nginx configuration file /etc/nginx/http.d/default.conf not found!"
	exit 1
fi

# Ensure nginx can write to log directories
mkdir -p /var/log/nginx /var/run
touch /var/log/nginx/access.log /var/log/nginx/error.log
chmod 755 /var/log/nginx /var/run
chmod 644 /var/log/nginx/*.log 2>/dev/null || true

# Verify nginx binary exists and is executable
if ! command -v nginx >/dev/null 2>&1; then
	bashio::log.error "Nginx binary not found!"
	exit 1
fi

# Final configuration test
if ! nginx -t; then
	bashio::log.error "Nginx configuration test failed before startup!"
	exit 1
fi

# Start nginx in background first to verify it starts successfully
bashio::log.info "Starting Nginx process (foreground mode)..."
nginx -g "daemon on;" || {
	bashio::log.error "Failed to start Nginx!"
	exit 1
}

# Wait a moment for nginx to start
sleep 2

# Verify nginx is actually running and listening
if ! pgrep -f "nginx: master process" >/dev/null 2>&1; then
	bashio::log.error "Nginx process not found after startup!"
	exit 1
fi

# Check if nginx is listening on port 80
if command -v ss >/dev/null 2>&1; then
	if ss -tuln | grep -q ":80 "; then
		bashio::log.info "✓ Nginx is listening on port 80"
		# Show what's listening
		ss -tuln | grep ":80 " || true
	else
		bashio::log.error "✗ Nginx is NOT listening on port 80!"
		bashio::log.error "Port status:"
		ss -tuln | grep ":80 " || bashio::log.error "No process found on port 80"
		# Try to get nginx error log
		if [ -f /var/log/nginx/error.log ]; then
			bashio::log.error "Last 10 lines of nginx error log:"
			tail -10 /var/log/nginx/error.log || true
		fi
		exit 1
	fi
elif command -v netstat >/dev/null 2>&1; then
	if netstat -tuln | grep -q ":80 "; then
		bashio::log.info "✓ Nginx is listening on port 80"
	else
		bashio::log.error "✗ Nginx is NOT listening on port 80!"
		exit 1
	fi
fi

# Test if nginx responds to requests
bashio::log.info "Testing nginx response..."
if command -v curl >/dev/null 2>&1; then
	RESPONSE=$(curl -s -o /dev/null -w "%{http_code}" http://127.0.0.1:80/ 2>/dev/null || echo "000")
	if [ "$RESPONSE" = "200" ] || [ "$RESPONSE" = "302" ] || [ "$RESPONSE" = "301" ]; then
		bashio::log.info "✓ Nginx responds to requests (HTTP $RESPONSE)"
	else
		bashio::log.warning "⚠ Nginx responded with HTTP $RESPONSE (expected 200/302/301)"
	fi
fi

# Check PHP-FPM connectivity
bashio::log.info "Testing PHP-FPM connectivity..."
if command -v nc >/dev/null 2>&1; then
	if nc -z 127.0.0.1 9000 2>/dev/null; then
		bashio::log.info "✓ PHP-FPM is reachable on 127.0.0.1:9000"
	else
		bashio::log.warning "⚠ PHP-FPM is not reachable on 127.0.0.1:9000"
	fi
fi

# Show final status
bashio::log.info "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
bashio::log.info "Nginx Status:"
bashio::log.info "  ✓ Process: Running"
bashio::log.info "  ✓ Listening: Port 80 (external: 8099)"
bashio::log.info "  ✓ Document Root: /var/www/wordpress"
bashio::log.info ""
bashio::log.info "Access WordPress at: http://<your-home-assistant-ip>:8099"
bashio::log.info "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"

# Restart nginx in foreground mode for proper signal handling
nginx -s stop 2>/dev/null || true
sleep 1

# Start nginx in foreground (for s6-overlay)
exec nginx -g "daemon off;"
