# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/svc-redis/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "Starting Redis..."
exec redis-server --save "" --appendonly no
