# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/svc-redis/run
#!/usr/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Redis..."
exec redis-server --save "" --appendonly no
