# netboot-xyz/root/etc/services.d/tftp/run
#!/usr/bin/with-contenv bash

/usr/sbin/in.tftpd \
	--foreground --listen --secure ${PORT_RANGE:+--port-range $PORT_RANGE} /config/menus
