# netboot-xyz/root/etc/services.d/tftp/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# TFTP must always use port 69 for PXE boot to work
# This is fixed and cannot be changed due to host_network mode
bashio::log.info "Starting TFTP server on port 69 (fixed for host_network mode and PXE boot)"

# Ensure PORT_RANGE is not set (would override default port 69)
unset PORT_RANGE

# Start TFTP server on port 69 (default port, required for PXE boot)
# --listen: listen on all interfaces
# --secure: only allow files within the specified directory
exec /usr/sbin/in.tftpd \
	--foreground --listen --secure /config/menus
