# planka/rootfs/etc/s6-overlay/s6-rc.d/svc-postgresql/run
#!/usr/bin/env bashio

bashio::log.info "Starting PostgreSQL..."

if ! bashio::fs.directory_exists "/run/postgresql"; then
    mkdir -p /run/postgresql
    chown postgres:postgres /run/postgresql
fi

exec su-exec postgres postgres -D /data/planka/database
