# antigravity-server/rootfs/etc/s6-overlay/s6-rc.d/init-antigravity/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: Antigravity-Server
# Initialize VNC password and environment
# ==============================================================================

bashio::log.info "Initializing Antigravity-Server..."

# Create required directories
mkdir -p /root/.config/tigervnc
mkdir -p /tmp/.X11-unix
chmod 1777 /tmp/.X11-unix

# Setup VNC password
if bashio::config.has_value 'vnc_password' && [ -n "$(bashio::config 'vnc_password')" ]; then
    bashio::log.info "Using configured VNC password from App options."
    VNC_PASSWORD=$(bashio::config 'vnc_password')
    if [ ${#VNC_PASSWORD} -gt 8 ]; then
        bashio::log.warning "User provided VNC password is longer than 8 characters. Truncating to the first 8 characters due to VNC protocol limitations."
        VNC_PASSWORD=$(echo "$VNC_PASSWORD" | head -c 8)
    fi
else
    # Ensure namespaced directory exists
    mkdir -p /data/antigravity_server

    if [ -f "/data/antigravity_server/vnc_password" ]; then
        bashio::log.info "Found existing VNC password in /data/antigravity_server/vnc_password."
        VNC_PASSWORD=$(cat "/data/antigravity_server/vnc_password")
        if [ -z "$VNC_PASSWORD" ]; then
             bashio::log.warning "Existing password file was empty. Generating new one."
             rm -f /data/antigravity_server/vnc_password
        else
             bashio::log.info "Using persisted random VNC password."
        fi
    fi

    if [ ! -f "/data/antigravity_server/vnc_password" ]; then
        bashio::log.info "Generating new random VNC password..."
        VNC_PASSWORD=$(python3 -c 'import secrets,string; print("".join(secrets.choice(string.ascii_letters + string.digits) for i in range(8)))')

        # atomic write with strict permissions
        echo "${VNC_PASSWORD}" > /data/antigravity_server/vnc_password.tmp
        chmod 600 /data/antigravity_server/vnc_password.tmp
        mv /data/antigravity_server/vnc_password.tmp /data/antigravity_server/vnc_password

        if [ -f "/data/antigravity_server/vnc_password" ]; then
            bashio::log.info "Successfully generated and persisted new random VNC password."
            bashio::log.info "---------------------------------------------------"
            bashio::log.info " VNC PASSWORD: ${VNC_PASSWORD}"
            bashio::log.info "---------------------------------------------------"
        else
            bashio::log.error "Failed to persist VNC password!"
        fi
    else
        bashio::log.info "VNC password already set. Find it in: /data/antigravity_server/vnc_password"
    fi
fi

# Set VNC password
touch /root/.config/tigervnc/passwd
chmod 600 /root/.config/tigervnc/passwd

bashio::log.info "Setting password using tigervncpasswd..."

if [ -x /usr/bin/tigervncpasswd ]; then
    VNCPASSWD_BIN="/usr/bin/tigervncpasswd"
else
    VNCPASSWD_BIN="vncpasswd"
fi

# Use filter mode (-f) to set password non-interactively
echo "${VNC_PASSWORD}" | $VNCPASSWD_BIN -f > /root/.config/tigervnc/passwd

if [ ${PIPESTATUS[1]} -ne 0 ]; then
    bashio::log.error "Failed to set VNC password"
    exit 1
fi

chmod 600 /root/.config/tigervnc/passwd
bashio::log.info "VNC password set successfully."

# Create xstartup script for XFCE4
# Also create Antigravity data directory (namespaced)
mkdir -p /data/antigravity_server/antigravity

mkdir -p /root/.config/autostart

# 1. Restore standard xstartup (always correct for VNC)
cat > /root/.config/tigervnc/xstartup << 'EOF'
#!/bin/bash
unset SESSION_MANAGER
unset DBUS_SESSION_BUS_ADDRESS
export XDG_CONFIG_DIRS=/etc/xdg:/etc/xfce4
# Start clipboard sync (Primary -> Clipboard and vice versa)
autocutsel -fork
autocutsel -selection PRIMARY -fork
# Bridge VNC clipboard to X11
vncconfig -nowin &
exec dbus-launch --exit-with-session startxfce4
EOF

chmod +x /root/.config/tigervnc/xstartup

# Configure Log Level
ANTIGRAVITY_LOG_LEVEL="info"
if bashio::config.has_value 'log_level'; then
    LOG_LEVEL=$(bashio::config 'log_level')
    case "${LOG_LEVEL}" in
        trace)       ANTIGRAVITY_LOG_LEVEL="trace" ;;
        debug)       ANTIGRAVITY_LOG_LEVEL="debug" ;;
        info|notice) ANTIGRAVITY_LOG_LEVEL="info" ;;
        warning)     ANTIGRAVITY_LOG_LEVEL="warn" ;;
        error)       ANTIGRAVITY_LOG_LEVEL="error" ;;
        fatal)       ANTIGRAVITY_LOG_LEVEL="critical" ;;
        *)           ANTIGRAVITY_LOG_LEVEL="info" ;;
    esac
fi

# 2. Configure Antigravity Autostart (Idempotent)
if bashio::config.true 'autostart_antigravity'; then
    # If the user (or upstream image) has already defined an autostart entry, we do NOT touch it.
    if [ ! -f "/root/.config/autostart/antigravity.desktop" ]; then
        bashio::log.info "Configuring default Antigravity autostart..."
        cat > /root/.config/autostart/antigravity.desktop << EOF
[Desktop Entry]
Type=Application
Name=Antigravity
# --no-sandbox is required because the desktop session runs as root.
# Electron applications (like VS Code/Antigravity) refuse to run as root without this flag.
Exec=antigravity --kiosk --no-sandbox --test-type --disable-dev-shm-usage --user-data-dir=/data/antigravity_server/antigravity --log ${ANTIGRAVITY_LOG_LEVEL}
StartupNotify=false
Terminal=false
Hidden=false
EOF
        chmod +x /root/.config/autostart/antigravity.desktop
    else
        # Update log level in existing file if needed (simple sed replacement if standard format)
        sed -i "s/--log [a-z]*/--log ${ANTIGRAVITY_LOG_LEVEL}/" /root/.config/autostart/antigravity.desktop || true
        bashio::log.info "Antigravity autostart entry already exists. Scaling back (idempotent)."
    fi
else
    bashio::log.info "Antigravity autostart is DISABLED by configuration."
    if [ -f "/root/.config/autostart/antigravity.desktop" ]; then
        bashio::log.info "Removing existing Antigravity autostart entry..."
        rm -f /root/.config/autostart/antigravity.desktop
    fi
fi


# 3. Set Google Chrome as default browser
bashio::log.info "Setting Google Chrome as default browser..."
mkdir -p /root/.config/google-chrome
# Suppress Chrome first run and welcome screens
touch "/root/.config/google-chrome/First Run"

# We manually create mimeapps.list because xdg-settings might fail without active X session
mkdir -p /root/.config
cat >> /root/.config/mimeapps.list <<EOF
[Default Applications]
text/html=google-chrome.desktop
x-scheme-handler/http=google-chrome.desktop
x-scheme-handler/https=google-chrome.desktop
x-scheme-handler/about=google-chrome.desktop
x-scheme-handler/unknown=google-chrome.desktop
EOF


# 4. Verify Tools
bashio::log.info "Verifying installed tools..."
if command -v git &> /dev/null; then
    bashio::log.info "Git version: $(git --version)"
else
    bashio::log.error "Git NOT found. Please check Dockerfile."
fi

if command -v docker &> /dev/null; then
    bashio::log.info "Docker version: $(docker --version)"
else
    bashio::log.error "Docker NOT found. Please check Dockerfile."
fi

bashio::log.info "VNC environment initialized"
