# antigravity-server/rootfs/etc/s6-overlay/s6-rc.d/init-antigravity/up
#!/bin/bash
/etc/s6-overlay/s6-rc.d/init-antigravity/run
