# wordpress/rootfs/etc/s6-overlay/s6-rc.d/svc-php-fpm/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Wait a moment for init to complete
sleep 1

# Verify PHP-FPM configuration
if [ ! -f /etc/php84/php-fpm.d/www.conf ]; then
	bashio::log.error "PHP-FPM configuration file not found!"
	exit 1
fi

bashio::log.info "Starting PHP-FPM..."
bashio::log.info "PHP-FPM will listen on 127.0.0.1:9000"

# Verify WordPress directory exists
if [ ! -d /var/www/wordpress ]; then
	bashio::log.error "WordPress directory /var/www/wordpress does not exist!"
	exit 1
fi

# Test PHP-FPM configuration before starting
if ! php-fpm84 -t 2>/dev/null; then
	bashio::log.warning "PHP-FPM configuration test had warnings, but continuing..."
fi

# Start PHP-FPM in foreground
bashio::log.info "Starting PHP-FPM process..."
bashio::log.info "PHP-FPM will be accessible at 127.0.0.1:9000 for Nginx"

# Start PHP-FPM and capture any errors
exec php-fpm84 -F || {
	bashio::log.error "PHP-FPM failed to start!"
	bashio::log.error "Check PHP-FPM configuration and logs for details."
	exit 1
}
