# netboot-xyz/root/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bash
exec \
/usr/sbin/dnsmasq --conf-dir=/etc/dnsmasq.d --no-daemon
