# vaultwarden/rootfs/etc/s6-overlay/s6-rc.d/vaultwarden/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Vaultwarden
# Vaultwarden main service
# ==============================================================================

# Bind to localhost as Nginx proxies to it
export ROCKET_ADDRESS=127.0.0.1
export ROCKET_PORT=7277
export ROCKET_DATA_DIR=/data

# Check and set Request Size Limit
if bashio::config.has_value 'request_size_limit'; then
    export ROCKET_LIMITS="{json=$(bashio::config 'request_size_limit')}"
fi

bashio::log.info "Starting Vaultwarden..."
exec /opt/vaultwarden/vaultwarden
