# vaultwarden/rootfs/etc/s6-overlay/s6-rc.d/vaultwarden/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Vaultwarden
# Runs the Vaultwarden server
# ==============================================================================

# Set environment from persistent files
export DATA_FOLDER=/data
export ROCKET_PORT=80
export ROCKET_WORKERS=2

if [ -f "/data/rocket_secret" ]; then
    export ROCKET_SECRET_KEY=$(cat /data/rocket_secret)
fi

if [ -f "/data/admin_token" ]; then
    export ADMIN_TOKEN=$(cat /data/admin_token)
fi

# API request size limit
if bashio::config.has_value 'request_size_limit'; then
    request_size_limit=$(bashio::config 'request_size_limit')
    export ROCKET_LIMITS="{json=${request_size_limit}}"
fi

# Run the Bitwarden server
bashio::log.info 'Starting the Vaultwarden server...'
cd /opt || bashio::exit.nok
exec ./vaultwarden
