# ShieldFile/rootfs/etc/services.d/filebrowser/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC1091

# Retrieve variables passed from init script or config
DB_PATH="/data/database.db"
BASE_DIR=$(bashio::config 'base_directory')

# Read Ingress Entry from file created by init script
if [ -f /var/run/shieldfile/ingress_entry ]; then
    INGRESS_ENTRY=$(cat /var/run/shieldfile/ingress_entry)
else
    bashio::log.warning "Ingress Entry file not found, defaulting to /"
    INGRESS_ENTRY="/"
fi

bashio::log.info "🚀 Starting FileBrowser backend on 127.0.0.1:8080..."
exec filebrowser \
    --address "127.0.0.1" \
    --port "8080" \
    --database "$DB_PATH" \
    --root "$BASE_DIR" \
    --baseurl "$INGRESS_ENTRY" \
    --log stdout
