# wiki.js3/rootfs/etc/services.d/wikijs/run
#!/usr/bin/with-contenv bashio

password=$(bashio::config 'db_password')

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL to be ready..."
for i in {1..30}; do
	if su-exec postgres pg_isready -q; then
		break
	fi
	sleep 1
done

if ! su-exec postgres pg_isready -q; then
	bashio::log.error "PostgreSQL failed to start!"
	exit 1
fi

bashio::log.info "PostgreSQL is ready!"

# Set postgres password if not already set
su-exec postgres psql -c "ALTER USER postgres WITH PASSWORD '${password}';" 2>/dev/null || true

#Drop database based on config flag (with safeguards)
if bashio::config.true 'reset_database'; then
	if ! bashio::config.true 'reset_database_confirm'; then
		bashio::log.error "Database reset requires both 'reset_database' and 'reset_database_confirm' to be enabled!"
		bashio::log.error "This is a DESTRUCTIVE operation that will DELETE ALL DATA!"
		bashio::log.error "Set 'reset_database_confirm: true' in your configuration to proceed."
		# Don't exit, just don't reset
	else
		bashio::log.warning "=========================================="
		bashio::log.warning "WARNING: DESTRUCTIVE DATABASE OPERATION"
		bashio::log.warning "=========================================="
		bashio::log.warning "This will DROP the entire 'wiki' database!"
		bashio::log.warning "All data will be permanently lost!"
		bashio::log.warning "=========================================="

		# Create timestamped backup before dropping
		BACKUP_DIR="/config/backups"
		mkdir -p "$BACKUP_DIR"
		TIMESTAMP=$(date +%Y%m%d_%H%M%S)
		BACKUP_FILE="$BACKUP_DIR/wiki_backup_${TIMESTAMP}.sql"

		bashio::log.info "Creating backup before database reset..."
		bashio::log.info "Backup location: $BACKUP_FILE"

		# Create backup using pg_dump (using local socket connection)
		if su-exec postgres pg_dump wiki >"$BACKUP_FILE" 2>/dev/null; then
			bashio::log.info "Backup created successfully: $BACKUP_FILE"
		else
			bashio::log.warning "Backup failed (database may not exist yet), continuing with reset..."
		fi

		bashio::log.warning 'Recreating database (dropping existing if present)...'
		su-exec postgres psql -d postgres -c "DROP DATABASE IF EXISTS wiki;"

		#Remove reset_database options
		bashio::addon.option 'reset_database'
		bashio::addon.option 'reset_database_confirm'
	fi
fi

# Create database if not exists
su-exec postgres psql -d postgres -tc "SELECT 1 FROM pg_database WHERE datname = 'wiki'" | grep -q 1 || su-exec postgres psql -d postgres -c "CREATE DATABASE wiki;"

# Enable pg_trgm extension (required for Wiki.js 3.0 search)
# Fixed: Install into pg_catalog to ensure visibility across all search_paths.
# Also DROP before CREATE to handle restarts where it might already exist.
bashio::log.info "Ensuring 'pg_trgm' extension is enabled..."
su-exec postgres psql -d wiki -c "CREATE EXTENSION IF NOT EXISTS pg_trgm SCHEMA pg_catalog;"

# Patch servers.mjs to avoid ERR_SERVER_ALREADY_LISTEN
bashio::log.info "Patching servers.mjs to fix bootloop..."
node -e "
const fs = require('fs');
const path = '/wiki/server/core/servers.mjs';
try {
  if (fs.existsSync(path)) {
    let content = fs.readFileSync(path, 'utf8');
    const regex = /((?:this\.)?\w+)\.listen\(/g;
    const newContent = content.replace(regex, 'if(!\$1.listening)\$1.listen(');
    if (content !== newContent) {
      fs.writeFileSync(path, newContent);
      console.log('Successfully patched servers.mjs');
    } else {
      console.log('Pattern not found or already patched in servers.mjs');
    }
  } else {
    console.log('servers.mjs not found, skipping patch during this run');
  }
} catch (e) {
  console.error('Error patching servers.mjs:', e);
}
"

echo "Starting Wiki.JS V3"
cd /wiki
# Configure Log Level
if bashio::config.has_value 'log_level'; then
    LOG_LEVEL=$(bashio::config 'log_level')
    case "${LOG_LEVEL}" in
        trace|debug) WIKI_LOG_LEVEL="debug" ;;
        info|notice) WIKI_LOG_LEVEL="info" ;;
        warning)     WIKI_LOG_LEVEL="warn" ;;
        error|fatal) WIKI_LOG_LEVEL="error" ;;
        *)           WIKI_LOG_LEVEL="info" ;;
    esac
    export LOG_LEVEL="$WIKI_LOG_LEVEL"
    bashio::log.info "Setting LOG_LEVEL to $WIKI_LOG_LEVEL"
fi

exec node server
