# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/init-ttrss-db/run
#!/usr/bin/env bashio

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL to be ready..."
max_retries=30
count=0
until pg_isready -h localhost -U postgres; do
    if [ $count -ge $max_retries ]; then
        bashio::exit.nok "PostgreSQL not ready after ${max_retries} attempts."
    fi
    bashio::log.info "Waiting for PostgreSQL..."
    sleep 1
    count=$((count + 1))
done

# Check password (created by init-tt-rss)
if bashio::fs.file_exists "/data/ttrss_db_password"; then
    TTRSS_DB_PASSWORD=$(cat /data/ttrss_db_password)
else
    bashio::exit.nok "Database password not found. Ensure init-tt-rss ran successfully."
fi

# Wait for config.php (generated by init-tt-rss)
max_retries=30
count=0
while [ ! -f /var/www/tt-rss/config.php ]; do
    if [ $count -ge $max_retries ]; then
        bashio::exit.nok "Timeout waiting for config.php generated by init-tt-rss"
    fi
    bashio::log.info "Waiting for config.php..."
    sleep 1
    count=$((count + 1))
done

# Securely Create Database and User
bashio::log.info "Ensuring TT-RSS database and user exist..."

# Create DB if missing
if [ "$(psql -h localhost -U postgres -v ON_ERROR_STOP=1 -tAc "SELECT 1 FROM pg_database WHERE datname = 'ttrss'")" != "1" ]; then
    psql -h localhost -U postgres -v ON_ERROR_STOP=1 -c "CREATE DATABASE ttrss;"
fi

# Create User if missing and Grant Privileges
# Use idempotent DO block to create user
# Note: DO block cannot run CREATE USER directly inside because it's a utility command.
# We stick to the check-then-create pattern but use safe variable passing.
if [ "$(psql -h localhost -U postgres -v ON_ERROR_STOP=1 -tAc "SELECT 1 FROM pg_roles WHERE rolname = 'ttrss'")" != "1" ]; then
    psql -h localhost -U postgres -v ON_ERROR_STOP=1 -v pw="'${TTRSS_DB_PASSWORD}'" -c "CREATE USER ttrss WITH PASSWORD :pw;"
fi

psql -h localhost -U postgres -v ON_ERROR_STOP=1 -c "GRANT ALL PRIVILEGES ON DATABASE ttrss TO ttrss;"

# Schema Update
bashio::log.info "Updating TT-RSS schema..."
cd /var/www/tt-rss
if ! s6-setuidgid nginx php ./update.php --update-schema=force-yes; then
    bashio::log.warning "Schema update returned error. This might be normal if already up to date, or fatal."
fi
