# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/init-ttrss-db/run
#!/usr/bin/with-contenv bashio

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL to be ready..."
max_retries=30
count=0
until pg_isready -U postgres; do
    if [ $count -ge $max_retries ]; then
        bashio::exit.nok "PostgreSQL not ready after ${max_retries} attempts."
    fi
    bashio::log.info "Waiting for PostgreSQL..."
    sleep 1
    count=$((count + 1))
done

# Check password (created by init-tt-rss)
if bashio::fs.file_exists "/data/ttrss_db_password"; then
    TTRSS_DB_PASSWORD=$(cat /data/ttrss_db_password)
else
    bashio::exit.nok "Database password not found. Ensure init-tt-rss ran successfully."
fi

# Wait for config.php (generated by init-tt-rss)
max_retries=30
count=0
while [ ! -f /var/www/tt-rss/config.php ]; do
    if [ $count -ge $max_retries ]; then
        bashio::exit.nok "Timeout waiting for config.php generated by init-tt-rss"
    fi
    bashio::log.info "Waiting for config.php..."
    sleep 1
    count=$((count + 1))
done

# Database Reset Check (Two-Step)
if bashio::config.true 'reset_database'; then
	if ! bashio::config.true 'reset_database_confirm'; then
		bashio::log.error "Database reset requires both 'reset_database' and 'reset_database_confirm' to be enabled!"
		bashio::log.error "This is a DESTRUCTIVE operation that will DELETE ALL DATA!"
		bashio::log.error "Set 'reset_database_confirm: true' in your configuration to proceed."
		bashio::exit.nok "Database reset confirmation required"
	fi

	bashio::log.warning "=========================================="
	bashio::log.warning "   ⚠️  DATABASE RESET ENABLED  ⚠️"
	bashio::log.warning "=========================================="
	bashio::log.warning "ALL TT-RSS DATA WILL BE PERMANENTLY DELETED!"
	bashio::log.warning "This includes:"
	bashio::log.warning "  - All feeds and subscriptions"
	bashio::log.warning "  - All articles and starred items"
	bashio::log.warning "  - All user preferences and filters"
	bashio::log.warning "=========================================="

	# Create timestamped backup before deleting
	BACKUP_DIR="/data/backups"
	mkdir -p "$BACKUP_DIR"
	TIMESTAMP=$(date +%Y%m%d_%H%M%S)
	BACKUP_FILE="$BACKUP_DIR/ttrss_backup_${TIMESTAMP}.sql"

	bashio::log.info "Creating backup before database reset..."
	bashio::log.info "Backup location: $BACKUP_FILE"

	# Create backup using pg_dump
	if psql -U postgres -lqt | cut -d \| -f 1 | grep -qw ttrss; then
		if pg_dump -U postgres ttrss > "$BACKUP_FILE" 2>/dev/null; then
			bashio::log.info "Backup created successfully: $BACKUP_FILE"
		else
			bashio::log.warning "Backup failed, continuing with reset..."
		fi
	else
		bashio::log.info "Database does not exist yet, skipping backup..."
	fi

	# Wait 3 seconds to give user time to cancel
	bashio::log.warning "Starting reset in 3 seconds... (Stop App NOW to abort)"
	sleep 3

	bashio::log.info "Proceeding with database reset..."

	# Drop database if it exists
	psql -U postgres -c "DROP DATABASE IF EXISTS ttrss;" || bashio::log.warning "Failed to drop database"

	bashio::log.info "=========================================="
	bashio::log.info "   ✅ DATABASE RESET COMPLETE"
	bashio::log.info "=========================================="
	bashio::log.info "Database has been dropped."
	bashio::log.info "It will be recreated on next step."
	bashio::log.warning "IMPORTANT: Disable 'reset_database' and 'reset_database_confirm' in the App settings!"
	bashio::log.warning "Otherwise, the database will be wiped again on next restart."
	bashio::log.info "=========================================="

	# Remove reset_database options
	bashio::app.option 'reset_database'
	bashio::app.option 'reset_database_confirm'
fi

# Securely Create Database and User
bashio::log.info "Ensuring TT-RSS database and user exist..."

# Create DB if missing
if [ "$(psql -U postgres -v ON_ERROR_STOP=1 -tAc "SELECT 1 FROM pg_database WHERE datname = 'ttrss'")" != "1" ]; then
    psql -U postgres -v ON_ERROR_STOP=1 -c "CREATE DATABASE ttrss;"
fi

# Create User if missing and Grant Privileges
# Use psql with escaped password (avoid embedding quotes in -v)
USER_EXISTS=$(psql -U postgres -v ON_ERROR_STOP=1 -tAc "SELECT 1 FROM pg_roles WHERE rolname = 'ttrss'" 2>/dev/null || echo "0")
if [ "$USER_EXISTS" != "1" ]; then
    bashio::log.info "Creating ttrss database user..."
    psql -U postgres -v ON_ERROR_STOP=1 -c "CREATE USER ttrss WITH PASSWORD '${TTRSS_DB_PASSWORD//\'/\'\'}'"
else
    # User exists - update password to ensure it matches the password file
    bashio::log.info "Updating ttrss database user password..."
    psql -U postgres -v ON_ERROR_STOP=1 -c "ALTER USER ttrss WITH PASSWORD '${TTRSS_DB_PASSWORD//\'/\'\'}'"
fi

psql -U postgres -v ON_ERROR_STOP=1 -c "GRANT ALL PRIVILEGES ON DATABASE ttrss TO ttrss;"

# PostgreSQL 15+ requires explicit schema permissions
# Grant permissions on public schema and make ttrss the owner of objects
psql -U postgres -d ttrss -v ON_ERROR_STOP=1 -c "GRANT ALL ON SCHEMA public TO ttrss;"
psql -U postgres -d ttrss -v ON_ERROR_STOP=1 -c "ALTER SCHEMA public OWNER TO ttrss;"

# Schema Update
bashio::log.info "Updating TT-RSS schema..."
cd /var/www/tt-rss
if ! s6-setuidgid nginx php ./update.php --update-schema=force-yes; then
    bashio::log.warning "Schema update returned error. This might be normal if already up to date, or fatal."
fi
