# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/init-ttrss-db/run
#!/usr/bin/env bashio

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL to be ready..."
max_retries=30
count=0
until pg_isready -U postgres; do
    if [ $count -ge $max_retries ]; then
        bashio::exit.nok "PostgreSQL not ready after ${max_retries} attempts."
    fi
    bashio::log.info "Waiting for PostgreSQL..."
    sleep 1
    count=$((count + 1))
done

# Check password (created by init-tt-rss)
if bashio::fs.file_exists "/data/ttrss_db_password"; then
    TTRSS_DB_PASSWORD=$(cat /data/ttrss_db_password)
else
    bashio::exit.nok "Database password not found. Ensure init-tt-rss ran successfully."
fi

# Wait for config.php (generated by init-tt-rss)
max_retries=30
count=0
while [ ! -f /var/www/tt-rss/config.php ]; do
    if [ $count -ge $max_retries ]; then
        bashio::exit.nok "Timeout waiting for config.php generated by init-tt-rss"
    fi
    bashio::log.info "Waiting for config.php..."
    sleep 1
    count=$((count + 1))
done

# Securely Create Database and User
bashio::log.info "Ensuring TT-RSS database and user exist..."

# Create DB if missing
if [ "$(psql -U postgres -v ON_ERROR_STOP=1 -tAc "SELECT 1 FROM pg_database WHERE datname = 'ttrss'")" != "1" ]; then
    psql -U postgres -v ON_ERROR_STOP=1 -c "CREATE DATABASE ttrss;"
fi

# Create User if missing and Grant Privileges
<<<<<<< Updated upstream
# Use psql with escaped password (avoid embedding quotes in -v)
USER_EXISTS=$(psql -U postgres -v ON_ERROR_STOP=1 -tAc "SELECT 1 FROM pg_roles WHERE rolname = 'ttrss'" 2>/dev/null || echo "0")
if [ "$USER_EXISTS" != "1" ]; then
    bashio::log.info "Creating ttrss database user..."
    psql -U postgres -v ON_ERROR_STOP=1 -c "CREATE USER ttrss WITH PASSWORD '${TTRSS_DB_PASSWORD//\'/\'\'}'"
else
    # User exists - update password to ensure it matches the password file
    bashio::log.info "Updating ttrss database user password..."
    psql -U postgres -v ON_ERROR_STOP=1 -c "ALTER USER ttrss WITH PASSWORD '${TTRSS_DB_PASSWORD//\'/\'\'}'"
fi

psql -U postgres -v ON_ERROR_STOP=1 -c "GRANT ALL PRIVILEGES ON DATABASE ttrss TO ttrss;"
=======
# Use psql with escaped password (avoid embedding quotes incorrectly)
USER_EXISTS=$(psql -h localhost -U postgres -v ON_ERROR_STOP=1 -tAc "SELECT 1 FROM pg_roles WHERE rolname = 'ttrss'" 2>/dev/null || echo "0")
if [ "$USER_EXISTS" != "1" ]; then
    bashio::log.info "Creating ttrss database user..."
    psql -h localhost -U postgres -v ON_ERROR_STOP=1 -c "CREATE USER ttrss WITH PASSWORD '${TTRSS_DB_PASSWORD//\'/\'\'}'"
else
    # User exists - update password to ensure it matches the password file
    bashio::log.info "Updating ttrss database user password..."
    psql -h localhost -U postgres -v ON_ERROR_STOP=1 -c "ALTER USER ttrss WITH PASSWORD '${TTRSS_DB_PASSWORD//\'/\'\'}'"
fi

# Grant privileges on database
psql -h localhost -U postgres -v ON_ERROR_STOP=1 -c "GRANT ALL PRIVILEGES ON DATABASE ttrss TO ttrss;"
>>>>>>> Stashed changes

# PostgreSQL 15+ requires explicit schema permissions
# Grant permissions on public schema and make ttrss the owner of objects
psql -h localhost -U postgres -d ttrss -v ON_ERROR_STOP=1 -c "GRANT ALL ON SCHEMA public TO ttrss;"
psql -h localhost -U postgres -d ttrss -v ON_ERROR_STOP=1 -c "ALTER SCHEMA public OWNER TO ttrss;"

# Schema Update
bashio::log.info "Updating TT-RSS schema..."
cd /var/www/tt-rss
if ! s6-setuidgid nginx php ./update.php --update-schema=force-yes; then
    bashio::log.warning "Schema update returned error. This might be normal if already up to date, or fatal."
fi
