# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/init-ttrss-db/up
#!/command/execlineb
foreground { /etc/s6-overlay/s6-rc.d/init-ttrss-db/run }
