# ShieldDNS/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Nginx..."
exec nginx -g "daemon off;"
