# komodo/rootfs/etc/s6-overlay/s6-rc.d/init-komodo/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

# Copy default config if needed or setup directories
mkdir -p /data/komodo
