# komodo/rootfs/etc/s6-overlay/s6-rc.d/init-komodo/up
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

# Setup directories
mkdir -p /data/komodo
