# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/init-backend/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Initializing Backend..."

bashio::net.wait_for 5432 localhost

cd /app

if [ -f "alembic.ini" ]; then
    bashio::log.info "Running database migrations..."

    export POSTGRES_SERVER="localhost"
    export POSTGRES_USER="postgres"
    export POSTGRES_DB="er_startseite"
    export POSTGRES_PASSWORD="${POSTGRES_PASSWORD}"

    # alembic upgrade head || bashio::log.warning "Migrations failed (optional)"
fi
