# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/init-backend/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Initializing Backend..."

bashio::net.wait_for 5432 localhost 30

cd /app

bashio::log.info "Backend initialization complete."
