# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting PHP-FPM..."
exec php-fpm84 -F
