# dbf/rootfs/etc/s6-overlay/s6-rc.d/dbf-periodic/run
#!/usr/bin/with-contenv bashio
exec /usr/sbin/crond -f -l 8
