# nginx/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio

# Start PHP-FPM in foreground
# Check if php-fpm84 exists, otherwise try php-fpm
if command -v php-fpm84 >/dev/null 2>&1; then
    exec php-fpm84 -F
elif command -v php-fpm >/dev/null 2>&1; then
    exec php-fpm -F
else
    bashio::log.warning "PHP-FPM not found, PHP support may not work"
    # Sleep to prevent restart loop if optional
    sleep 365d
fi
