# sap-abap-cloud-dev/rootfs/etc/s6-overlay/s6-rc.d/sap-abap/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# SAP ABAP Cloud Developer Trial - Service Runner
# ==============================================================================

# System Requirements (from SAP documentation)
# Recommended: 4 CPUs, 16GB RAM, 150GB Disk
# Minimum (with ignore_requirements): 4 CPUs, 8GB RAM, 50GB Disk
REQUIRED_CPUS=4
RECOMMENDED_RAM_GB=16
MINIMUM_RAM_GB=8
RECOMMENDED_DISK_GB=150
MINIMUM_DISK_GB=50

bashio::log.info "=========================================="
bashio::log.info " SAP ABAP Cloud Developer Trial"
bashio::log.info " FaserF's Home Assistant Apps"
bashio::log.info "=========================================="
bashio::log.info ""

# ------------------------------------------------------------------------------
# License and Disclaimer
# ------------------------------------------------------------------------------
bashio::log.warning "=========================================="
bashio::log.warning " IMPORTANT DISCLAIMERS"
bashio::log.warning "=========================================="
bashio::log.warning ""
bashio::log.warning "NO LICENSE PROVIDED: This App does NOT"
bashio::log.warning "include any SAP license. You must obtain"
bashio::log.warning "your own license from SAP."
bashio::log.warning ""
bashio::log.warning "NO WARRANTY: This App is provided AS IS."
bashio::log.warning "The maintainer assumes NO LIABILITY for data"
bashio::log.warning "loss, system damage, or any other issues."
bashio::log.warning ""
bashio::log.warning "FOR TESTING ONLY: This App is intended"
bashio::log.warning "solely for personal learning and testing."
bashio::log.warning ""
bashio::log.warning "=========================================="
bashio::log.info ""

# ------------------------------------------------------------------------------
# License Agreement Check
# ------------------------------------------------------------------------------
if ! bashio::config.true 'agree_to_license'; then
    bashio::log.fatal "=========================================="
    bashio::log.fatal " LICENSE AGREEMENT REQUIRED"
    bashio::log.fatal "=========================================="
    bashio::log.fatal ""
    bashio::log.fatal "You must agree to SAP's license terms!"
    bashio::log.fatal "Set 'agree_to_license: true' in your"
    bashio::log.fatal "App configuration."
    bashio::log.fatal ""
    bashio::log.fatal "SAP License Terms:"
    bashio::log.fatal "https://www.sap.com/about/legal/disclaimer.html"
    bashio::log.fatal ""
    bashio::log.fatal "By setting agree_to_license: true, you confirm:"
    bashio::log.fatal "- You have read and accept SAP's license terms"
    bashio::log.fatal "- You understand this is for testing only"
    bashio::log.fatal "- You accept all liability for your usage"
    bashio::log.fatal "=========================================="
    bashio::exit.nok
fi
bashio::log.info "SAP License agreement accepted by user."

# ------------------------------------------------------------------------------
# System Requirements Check
# ------------------------------------------------------------------------------
bashio::log.info "Checking system requirements..."
if ! LOG_LEVEL=$(bashio::config 'log_level') || [ -z "$LOG_LEVEL" ]; then
    bashio::log.warning "Failed to fetch log_level configuration. Using default: info"
    LOG_LEVEL="info"
fi
export LOG_LEVEL
bashio::log.level "${LOG_LEVEL}"

IGNORE_REQUIREMENTS=$(bashio::config 'ignore_requirements')

# Get system information
CPU_COUNT=$(nproc 2>/dev/null || grep -c processor /proc/cpuinfo 2>/dev/null || echo "0")
TOTAL_RAM_KB=$(grep MemTotal /proc/meminfo 2>/dev/null | awk '{print $2}' || echo "0")
TOTAL_RAM_GB=$((TOTAL_RAM_KB / 1024 / 1024))
AVAILABLE_DISK_KB=$(df / 2>/dev/null | tail -1 | awk '{print $4}' || echo "0")
AVAILABLE_DISK_GB=$((AVAILABLE_DISK_KB / 1024 / 1024))

bashio::log.info "Detected system resources:"
bashio::log.info "  - CPUs: ${CPU_COUNT}"
bashio::log.info "  - RAM: ${TOTAL_RAM_GB} GB"
bashio::log.info "  - Available Disk: ${AVAILABLE_DISK_GB} GB"
bashio::log.info ""

# CPU Check (always required, no bypass)
if [[ "${CPU_COUNT}" -lt "${REQUIRED_CPUS}" ]]; then
    bashio::log.fatal "Insufficient CPUs: ${CPU_COUNT} detected, ${REQUIRED_CPUS} required."
    bashio::log.fatal "SAP ABAP Trial requires at least ${REQUIRED_CPUS} CPUs."
    bashio::exit.nok
fi

# RAM Check
if [[ "${IGNORE_REQUIREMENTS}" == "true" ]]; then
    if [[ "${TOTAL_RAM_GB}" -lt "${MINIMUM_RAM_GB}" ]]; then
        bashio::log.fatal "Insufficient RAM: ${TOTAL_RAM_GB} GB detected, ${MINIMUM_RAM_GB} GB minimum required."
        bashio::exit.nok
    elif [[ "${TOTAL_RAM_GB}" -lt "${RECOMMENDED_RAM_GB}" ]]; then
        bashio::log.warning "RAM below recommended: ${TOTAL_RAM_GB} GB (recommended: ${RECOMMENDED_RAM_GB} GB)"
        bashio::log.warning "Performance may be severely degraded."
    fi
else
    if [[ "${TOTAL_RAM_GB}" -lt "${RECOMMENDED_RAM_GB}" ]]; then
        bashio::log.fatal "Insufficient RAM: ${TOTAL_RAM_GB} GB detected, ${RECOMMENDED_RAM_GB} GB recommended."
        bashio::log.fatal "Set 'ignore_requirements: true' to run with minimum ${MINIMUM_RAM_GB} GB RAM."
        bashio::exit.nok
    fi
fi

# Disk Check
if [[ "${IGNORE_REQUIREMENTS}" == "true" ]]; then
    if [[ "${AVAILABLE_DISK_GB}" -lt "${MINIMUM_DISK_GB}" ]]; then
        bashio::log.fatal "Insufficient disk: ${AVAILABLE_DISK_GB} GB available, ${MINIMUM_DISK_GB} GB minimum required."
        bashio::exit.nok
    elif [[ "${AVAILABLE_DISK_GB}" -lt "${RECOMMENDED_DISK_GB}" ]]; then
        bashio::log.warning "Disk below recommended: ${AVAILABLE_DISK_GB} GB (recommended: ${RECOMMENDED_DISK_GB} GB)"
    fi
else
    if [[ "${AVAILABLE_DISK_GB}" -lt "${RECOMMENDED_DISK_GB}" ]]; then
        bashio::log.fatal "Insufficient disk: ${AVAILABLE_DISK_GB} GB available, ${RECOMMENDED_DISK_GB} GB recommended."
        bashio::log.fatal "Set 'ignore_requirements: true' to run with minimum ${MINIMUM_DISK_GB} GB disk."
        bashio::exit.nok
    fi
fi

bashio::log.info "System requirements check passed!"
bashio::log.info ""

# ------------------------------------------------------------------------------
# Startup Information
# ------------------------------------------------------------------------------
bashio::log.info "=========================================="
bashio::log.info " STARTING SAP ABAP PLATFORM"
bashio::log.info "=========================================="
bashio::log.info ""
bashio::log.info "Initial startup can take 5-10 minutes."
bashio::log.info ""
bashio::log.info "Connection Details:"
bashio::log.info "  - Application Server: Your HA IP"
bashio::log.info "  - Instance Number: 00"
bashio::log.info "  - System ID: A4H"
bashio::log.info "  - Client: 001"
bashio::log.info "  - User: DEVELOPER"
bashio::log.info ""
bashio::log.info "Web: https://localhost:8443/sap/bc/ui2/flp"
bashio::log.info "=========================================="
bashio::log.info ""

# ------------------------------------------------------------------------------
# Execute SAP Startup
# ------------------------------------------------------------------------------
exec /usr/local/bin/startup.sh
