# planka/rootfs/etc/s6-overlay/s6-rc.d/svc-planka/run
#!/usr/bin/env bashio

# Set strict error handling
set -o pipefail

bashio::log.info "Starting Planka..."

# Configure environment variables
# DB Password
if bashio::fs.file_exists "/data/planka_db_password"; then
    # Read without trailing newline
    PLANKA_DB_PASSWORD=$(tr -d '\n' < /data/planka_db_password)
else
    # Should not happen if init-planka ran
    PLANKA_DB_PASSWORD="planka"
fi

export PLANKA_DB_PASSWORD
ENCODED_PASS=$(python3 -c "import urllib.parse, os; print(urllib.parse.quote(os.environ.get('PLANKA_DB_PASSWORD', ''), safe=''))")

export DATABASE_URL="postgresql://planka:${ENCODED_PASS}@localhost:5432/planka"
export PORT=1337

# Secret Key
SECRET=$(bashio::config 'secret_key')
if [ -z "$SECRET" ]; then
    if ! bashio::fs.file_exists "/data/secret_key"; then
        bashio::log.warning "No secret_key provided. Generating a persistent one."
        openssl rand -hex 32 > /data/secret_key
    fi
    SECRET=$(cat /data/secret_key)
fi
export SECRET_KEY="$SECRET"
export TRUST_PROXY=1

# BASE_URL handling
if bashio::config.has_value 'base_url'; then
    export BASE_URL="$(bashio::config 'base_url')"
elif bashio::var.has_value "$(bashio::addon.ingress_entry)"; then
    export BASE_URL="http://localhost$(bashio::addon.ingress_entry)"
else
    export BASE_URL="http://localhost:1337"
fi
export SAILS_SECURITY_CORS_ALLOW_ORIGINS='*'
# Log level - schema now uses Winston-compatible levels
# Winston supports: error, warn, info, http, verbose, debug, silly
if bashio::config.has_value 'log_level'; then
    export LOG_LEVEL="$(bashio::config 'log_level')"
else
    export LOG_LEVEL="info"
fi

if ! cd /opt/planka; then
    bashio::exit.nok "Failed to change directory to /opt/planka"
fi

# Start Planka
exec ./start.sh
