# planka/rootfs/etc/s6-overlay/s6-rc.d/svc-planka/run
#!/usr/bin/env bashio
# shellcheck disable=SC1091
source /usr/lib/bashio/app.sh
source /usr/lib/bashio/banner.sh

# Set strict error handling
set -o pipefail

bashio::log.info "Starting Planka..."

# Configure environment variables
# DB Password
if bashio::fs.file_exists "/data/planka/db_password"; then
    # Read without trailing newline
    PLANKA_DB_PASSWORD=$(tr -d '\n' < /data/planka/db_password)
else
    # Should not happen if init-planka ran
    PLANKA_DB_PASSWORD="planka"
fi

export PLANKA_DB_PASSWORD
ENCODED_PASS=$(python3 -c "import urllib.parse, os; print(urllib.parse.quote(os.environ.get('PLANKA_DB_PASSWORD', ''), safe=''))")

export DATABASE_URL="postgresql://planka:${ENCODED_PASS}@localhost:5432/planka"
export PORT=1337
export NODE_ENV=production

# Secret Key
SECRET=$(bashio::config 'secret_key')
if [ -z "$SECRET" ]; then
    if ! bashio::fs.file_exists "/data/secret_key"; then
        bashio::log.warning "No secret_key provided. Generating a persistent one."
        openssl rand -hex 32 > /data/secret_key
    fi
    SECRET=$(cat /data/secret_key)
fi
export SECRET_KEY="$SECRET"
export TRUST_PROXY="true"
export TRUST_PROXY_BOOLEAN="true"
export TRUST_PROXY_VALUES="true"
# Note: Planka 2.x uses TRUST_PROXY="true" string for boolean behavior
export TRUST_PROXY_REAL="true"

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL to be ready..."
until pg_isready -h localhost -U postgres > /dev/null; do
    sleep 1
done

# Database Safety Check: Ensure the planka database and privileges exist
bashio::log.info "Ensuring Planka database exists..."
PLANKA_DATA_DIR="/data/planka"
PLANKA_DB_PASSWORD=$(cat "${PLANKA_DATA_DIR}/db_password")
export PGPASSWORD=$(cat "${PLANKA_DATA_DIR}/postgres_password")

psql -h localhost -U postgres -tc "SELECT 1 FROM pg_database WHERE datname = 'planka'" | grep -q 1 || \
    psql -h localhost -U postgres -c "CREATE DATABASE planka;"

if [ "$(psql -h localhost -U postgres -v ON_ERROR_STOP=1 -tAc "SELECT 1 FROM pg_roles WHERE rolname = 'planka'")" != "1" ]; then
    psql -h localhost -U postgres -v ON_ERROR_STOP=1 -c "CREATE USER planka WITH PASSWORD '${PLANKA_DB_PASSWORD}';"
fi

psql -h localhost -U postgres -v ON_ERROR_STOP=1 -c "GRANT ALL PRIVILEGES ON DATABASE planka TO planka;" > /dev/null 2>&1
psql -h localhost -U postgres -d planka -c "GRANT ALL ON SCHEMA public TO planka;" > /dev/null 2>&1
psql -h localhost -U postgres -d planka -c "GRANT ALL PRIVILEGES ON ALL TABLES IN SCHEMA public TO planka;" > /dev/null 2>&1
psql -h localhost -U postgres -d planka -c "GRANT ALL PRIVILEGES ON ALL SEQUENCES IN SCHEMA public TO planka;" > /dev/null 2>&1
psql -h localhost -U postgres -d planka -c "ALTER DEFAULT PRIVILEGES IN SCHEMA public GRANT ALL ON TABLES TO planka;" > /dev/null 2>&1
psql -h localhost -U postgres -d planka -c "ALTER DEFAULT PRIVILEGES IN SCHEMA public GRANT ALL ON SEQUENCES TO planka;" > /dev/null 2>&1

bashio::app.print_banner

# BASE_URL handling - Highly robust detection with diagnostics
DETECTED_BASE_URL=""

if bashio::var.has_value "$(bashio::config 'base_url')"; then
    # Aggressively strip any provided protocol and trailing slash
    CLEAN_BASE=$(echo "$(bashio::config 'base_url')" | sed -E 's|^https?://||' | sed 's|/$||')
    
    # Strictly re-assemble based on addon SSL setting
    if bashio::config.true 'ssl'; then
        DETECTED_BASE_URL="https://${CLEAN_BASE}"
    else
        DETECTED_BASE_URL="http://${CLEAN_BASE}"
    fi
elif bashio::var.has_value "$(bashio::app.ingress_entry)"; then
    DETECTED_BASE_URL="http://localhost$(bashio::app.ingress_entry)"
fi

# Final fallback if nothing was detected or if detection produced an empty string
if [ -z "$DETECTED_BASE_URL" ] || [ "$DETECTED_BASE_URL" == "null" ]; then
    DETECTED_BASE_URL="http://localhost:1337"
fi

export BASE_URL="$DETECTED_BASE_URL"
bashio::log.info "Configured BASE_URL: ${BASE_URL}"

# Secure Cookie Handling
if bashio::config.true 'ssl'; then
    bashio::log.info "SSL enabled, configuring secure cookies..."
    export SESSION_COOKIE_SECURE="true"
else
    export SESSION_COOKIE_SECURE="false"
fi

# Log level - schema now uses Winston-compatible levels
# Winston supports: error, warn, info, http, verbose, debug, silly
if ! App_LOG_LEVEL=$(bashio::config 'log_level') || [ -z "$App_LOG_LEVEL" ]; then
    bashio::log.warning "Failed to fetch log_level configuration. Using default: info"
    App_LOG_LEVEL="info"
fi
case "$App_LOG_LEVEL" in
    all|trace)   export LOG_LEVEL="silly" ;;
    debug)        export LOG_LEVEL="debug" ;;
    info|notice)  export LOG_LEVEL="info"  ;;
    warning|warn) export LOG_LEVEL="warn"  ;;
    error|fatal)  export LOG_LEVEL="error" ;;
    *)            export LOG_LEVEL="info"  ;;
esac

# 4. Configure Secure CORS Whitelist
# We starts with trusted local origins
CORS_ORIGINS="http://localhost:1337,http://127.0.0.1:1337,http://homeassistant:1337,http://homeassistant.local:1337"

# Add Home Assistant's internal IP for robust local access
LOCAL_IP=$(getent hosts homeassistant | awk '{ print $1 }' | head -n 1)
if [ -n "$LOCAL_IP" ]; then
    CORS_ORIGINS="${CORS_ORIGINS},http://${LOCAL_IP}:8123,https://${LOCAL_IP}:8123,http://${LOCAL_IP}:1337"
fi

# Add the re-assembled BASE_URL origin to CORS whitelist
if [ -n "$BASE_URL" ]; then
    CORS_ORIGINS="${CORS_ORIGINS},${BASE_URL}"
fi

# Apply restricted CORS to Sails
export sails_security__cors__allowOrigins="${CORS_ORIGINS}"
export sails_security__cors__allRoutes='true'
export sails_security__cors__allowCredentials='true'

export SAILS_LOG_LEVEL="$LOG_LEVEL"
export sails_log__level="$LOG_LEVEL"

if ! cd /opt/planka; then
    bashio::exit.nok "Failed to change directory to /opt/planka"
fi

# 1. Database Safety & Migration Check
# Ensures that even if init-planka was skipped, the DB is ready.
bashio::log.info "Performing final database readiness check..."
export PGPASSWORD=$(cat /data/planka/postgres_password)
psql -h localhost -U postgres -tc "SELECT 1 FROM pg_database WHERE datname = 'planka'" | grep -q 1 || \
    psql -h localhost -U postgres -c "CREATE DATABASE planka;"

# Ensure tables exist and admin user is seeded
# This is safe to run multiple times
./node_modules/.bin/knex migrate:latest --cwd db > /dev/null 2>&1
if ! bashio::fs.file_exists "/data/planka/admin_created"; then
    bashio::log.info "Seeding initial admin user..."
    export DEFAULT_ADMIN_EMAIL="admin@planka.local"
    export DEFAULT_ADMIN_PASSWORD=$(cat /data/planka/admin_password)
    export DEFAULT_ADMIN_NAME="Administrator"
    export DEFAULT_ADMIN_USERNAME="admin"
    ./node_modules/.bin/knex seed:run --cwd db --specific default.js > /dev/null 2>&1 && touch /data/planka/admin_created
fi

# 2. Patch config/log.js for strict logging control
if [ -f config/log.js ]; then
    sed -i "s/\/\/ level: 'info',/level: process.env.SAILS_LOG_LEVEL || 'info',/g" config/log.js
fi

# 3. Patch Planka config for production hardening
if [ -f config/env/production.js ]; then
    bashio::log.info "Hardening production configuration..."
    
    # Enable trustProxy and secure cookies
    sed -i "s/trustProxy: false,/trustProxy: true,/g" config/env/production.js
    if [ "$SESSION_COOKIE_SECURE" == "true" ]; then
        sed -i 's/\/\/ secure: true,/secure: true,/g' config/env/production.js
        sed -i 's/secure: false,/secure: true,/g' config/env/production.js
    fi

    # Use beforeConnect to allow all socket origins in production.
    sed -i "s/onlyAllowOrigins: origins,/beforeConnect: (handshake, cb) => cb(null, true),/g" config/env/production.js
fi

# 4. Patch Planka HTTP middleware to be Ingress-aware
if [ -f config/http.js ]; then
    bashio::log.info "Patching Planka HTTP middleware for robust Ingress routing..."
    cat << 'EOF' > config/http.js
const serveStatic = require('serve-static');
const sails = require('sails');

module.exports.http = {
  trustProxy: true,
  middleware: {
    poweredBy: false,

    www(req, res, next) {
      // Prioritize X-Ingress-Path header followed by URL detection
      const ingressHeader = req.headers['x-ingress-path'];
      if (ingressHeader) {
        req.url = req.url.replace(/^\/api\/hassio_ingress\/[^\/]+\/?/, '/');
      } else if (req.url.startsWith('/api/hassio_ingress/')) {
        req.url = req.url.replace(/^\/api\/hassio_ingress\/[^\/]+\/?/, '/');
      }
      
      const absolutePublicDir = '/opt/planka/public';
      const middleware = serveStatic(absolutePublicDir, {
        maxAge: sails.config.http.cache,
        immutable: req.url.startsWith('/assets/'),
      });

      return middleware(req, res, next);
    },
  },
};
EOF
fi

# 5. Patch Planka to support Ingress dynamically in the UI
if [ -f api/controllers/index.js ]; then
    bashio::log.info "Patching Planka for smart dynamic Ingress support..."
    cat << 'EOF' > api/controllers/index.js
module.exports = {
  exits: {
    success: {
      responseType: 'view',
      viewTemplatePath: 'index',
    },
  },
  fn() {
    // Detect Ingress based on header or URL
    const ingressHeader = this.req.headers['x-ingress-path'];
    if (ingressHeader) {
        return { basePath: ingressHeader.replace(/\/$/, '') };
    }

    const originalUrl = this.req.originalUrl || '';
    const match = originalUrl.match(/^\/api\/hassio_ingress\/[^\/]+/);
    const basePath = match ? match[0] : '';
    
    return {
      basePath: basePath,
    };
  },
};
EOF
fi

# 6. Disable CSRF to allow proxy-based logins
if [ -f config/security.js ]; then
    sed -i "s/\/\/ csrf: false,/csrf: false,/g" config/security.js
    sed -i "s/csrf: true,/csrf: false,/g" config/security.js
fi

# 7. Start Planka with a safe, high-performance output filter
# We use a shell-level filter to remove noisy framework warnings without touching the source code.
# This prevents all SyntaxErrors and ensures perfect stability.
FILTER_SIGS="MemoryStore is not|designed for a production environment|will not scale past a single process|When specifying multiple origins|setting should be an array|session.cookie.secure|cookie will be sent over non-TLS|If your app is behind a proxy|trustProxy|sailsjs.com|For more help:|  • "

if [[ "$LOG_LEVEL" =~ ^(debug|silly)$ ]]; then
    # In debug mode, show everything
    bashio::log.info "Lifting Planka (Debug mode)..."
    exec ./start.sh
else
    # In production/info mode, filter the noise silently with unbuffered Python
    # This is 100% compatible and ensures logs appear instantly.
    bashio::log.info "Lifting Planka..."
    exec ./start.sh 2>&1 | python3 -u -c "import sys,re; [sys.stdout.write(l) for l in sys.stdin if not re.search(r'$FILTER_SIGS', l)]"
fi
