# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/init-postgres/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Initializing PostgreSQL..."

# Create runtime directory for PostgreSQL lock files
mkdir -p /run/postgresql
chown postgres:postgres /run/postgresql

# Database Reset Check (Two-Step)
if bashio::config.true 'reset_database'; then
	if ! bashio::config.true 'reset_database_confirm'; then
		bashio::log.error "Database reset requires both 'reset_database' and 'reset_database_confirm' to be enabled!"
		bashio::log.error "This is a DESTRUCTIVE operation that will DELETE ALL DATA!"
		bashio::log.error "Set 'reset_database_confirm: true' in your configuration to proceed."
		exit 1
	fi

	bashio::log.warning "=========================================="
	bashio::log.warning "   ⚠️  DATABASE RESET ENABLED  ⚠️"
	bashio::log.warning "=========================================="
	bashio::log.warning "ALL ER-DASHBOARD DATA WILL BE PERMANENTLY DELETED!"
	bashio::log.warning "=========================================="

	# Backup directory (namespaced)
	BACKUP_DIR="/data/er_dashboard/backups"
	mkdir -p "$BACKUP_DIR"
	TIMESTAMP=$(date +%Y%m%d_%H%M%S)

	# Only backup if database directory exists and has data
	if [ -d "/data/er_dashboard/database" ] && [ -n "$(ls -A "/data/er_dashboard/database" 2>/dev/null)" ]; then
		BACKUP_FILE="$BACKUP_DIR/er_dashboard_backup_${TIMESTAMP}.tar.gz"
		bashio::log.info "Creating backup: $BACKUP_FILE"
		tar -czf "$BACKUP_FILE" -C /data/er_dashboard database 2>/dev/null || bashio::log.warning "Backup failed"
	fi

	bashio::log.warning "Starting reset in 3 seconds..."
	sleep 3

	# Remove database directory
	if [ -d "/data/er_dashboard/database" ]; then
		rm -rf "/data/er_dashboard/database"
		bashio::log.info "Database directory removed"
	fi

	bashio::log.info "✅ DATABASE RESET COMPLETE"
	bashio::log.warning "IMPORTANT: Disable reset options in add-on settings!"

	bashio::addon.option 'reset_database'
	bashio::addon.option 'reset_database_confirm'
fi

if [ ! -d "/data/er_dashboard/database" ] || [ -z "$(ls -A "/data/er_dashboard/database")" ]; then
    bashio::log.info "Creating new database cluster in /data/er_dashboard/database..."
    mkdir -p /data/er_dashboard/database
    chown -R postgres:postgres /data/er_dashboard/database
    chmod 700 /data/er_dashboard/database

    echo "${POSTGRES_PASSWORD}" > /tmp/pgpass
    chown postgres:postgres /tmp/pgpass

    if ! su-exec postgres initdb -D /data/er_dashboard/database --username=postgres --auth=md5 --pwfile=/tmp/pgpass 2>&1; then
        bashio::log.error "initdb failed! Command: su-exec postgres initdb -D /data/er_dashboard/database --username=postgres --auth=md5 --pwfile=/tmp/pgpass"
        rm -f /tmp/pgpass
        exit 1
    fi

    rm -f /tmp/pgpass
else
    bashio::log.info "Database cluster exists in /data/er_dashboard/database."
    # Ensure permissions are correct even if it exists
    chown -R postgres:postgres /data/er_dashboard/database
    chmod 700 /data/er_dashboard/database
fi
