# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/init-postgres/run
#!/usr/bin/with-contenv bashio

# Load custom bashio libraries
# shellcheck source=/dev/null
. /usr/lib/bashio/app.sh
# shellcheck shell=bash

bashio::log.info "Initializing PostgreSQL..."

# Create runtime directory for PostgreSQL lock files
mkdir -p /run/postgresql
chown postgres:postgres /run/postgresql

# Database Reset Check (Two-Step)
if bashio::config.true 'reset_database'; then
	if ! bashio::config.true 'reset_database_confirm'; then
		bashio::log.error "Database reset requires both 'reset_database' and 'reset_database_confirm' to be enabled!"
		bashio::log.error "This is a DESTRUCTIVE operation that will DELETE ALL DATA!"
		bashio::log.error "Set 'reset_database_confirm: true' in your configuration to proceed."
		exit 1
	fi

	bashio::log.warning "=========================================="
	bashio::log.warning "   ⚠️  DATABASE RESET ENABLED  ⚠️"
	bashio::log.warning "=========================================="
	bashio::log.warning "ALL ER-DASHBOARD DATA WILL BE PERMANENTLY DELETED!"
	bashio::log.warning "=========================================="

	# Backup directory (namespaced)
	BACKUP_DIR="/data/er_dashboard/backups"
	mkdir -p "$BACKUP_DIR"
	TIMESTAMP=$(date +%Y%m%d_%H%M%S)

	# Only backup if database directory exists and has data
	if [ -d "/data/er_dashboard/database" ] && [ -n "$(ls -A "/data/er_dashboard/database" 2>/dev/null)" ]; then
		BACKUP_FILE="$BACKUP_DIR/er_dashboard_backup_${TIMESTAMP}.tar.gz"
		bashio::log.info "Creating backup: $BACKUP_FILE"
		tar -czf "$BACKUP_FILE" -C /data/er_dashboard database 2>/dev/null || bashio::log.warning "Backup failed"
	fi

	bashio::log.warning "Starting reset in 3 seconds..."
	sleep 3

	# Remove database directory
	if [ -d "/data/er_dashboard/database" ]; then
		rm -rf "/data/er_dashboard/database"
		bashio::log.info "Database directory removed"
	fi

	bashio::log.info "✅ DATABASE RESET COMPLETE"
	bashio::log.warning "IMPORTANT: Disable reset options in App settings!"

	bashio::app.option 'reset_database'
	bashio::app.option 'reset_database_confirm'
fi

if [ ! -d "/data/er_dashboard/database" ] || [ -z "$(ls -A "/data/er_dashboard/database")" ]; then
    bashio::log.info "Creating new database cluster in /data/er_dashboard/database..."
    mkdir -p /data/er_dashboard/database
    chown -R postgres:postgres /data/er_dashboard/database
    chmod 700 /data/er_dashboard/database

    echo "${POSTGRES_PASSWORD}" > /tmp/pgpass
    chown postgres:postgres /tmp/pgpass

    if ! su-exec postgres initdb -D /data/er_dashboard/database --username=postgres --auth=md5 --pwfile=/tmp/pgpass 2>&1; then
        bashio::log.error "initdb failed! Command: su-exec postgres initdb -D /data/er_dashboard/database --username=postgres --auth=md5 --pwfile=/tmp/pgpass"
        rm -f /tmp/pgpass
        exit 1
    fi

    rm -f /tmp/pgpass
else
    bashio::log.info "Database cluster exists in /data/er_dashboard/database."
    # Ensure permissions are correct even if it exists
    chown -R postgres:postgres /data/er_dashboard/database
    chmod 700 /data/er_dashboard/database

    # Version check & upgrade logic
    if [ -f "/data/er_dashboard/database/PG_VERSION" ]; then
        DB_VERSION=$(cat "/data/er_dashboard/database/PG_VERSION")
        NEW_VERSION=$(postgres -V | awk '{print $NF}' | cut -d. -f1)

        if [ "$DB_VERSION" -lt "$NEW_VERSION" ]; then
            bashio::log.warning "Database version mismatch detected!"
            bashio::log.warning "Current database version: $DB_VERSION"
            bashio::log.warning "Target PostgreSQL version: $NEW_VERSION"
            bashio::log.info "Attempting automatic database upgrade using pg_upgrade..."

            OLD_BINDIR="/usr/libexec/postgresql${DB_VERSION}"
            NEW_BINDIR="/usr/libexec/postgresql${NEW_VERSION}"

            if [ ! -d "$NEW_BINDIR" ]; then
                if [ -f "/usr/libexec/postgresql/pg_upgrade" ]; then
                    NEW_BINDIR="/usr/libexec/postgresql"
                else
                    NEW_BINDIR="/usr/bin"
                fi
            fi

            if [ ! -d "$OLD_BINDIR" ]; then
                bashio::log.error "Old PostgreSQL binary directory $OLD_BINDIR not found! Cannot perform upgrade."
                exit 1
            fi

            bashio::log.info "Using old bindir: $OLD_BINDIR"
            bashio::log.info "Using new bindir: $NEW_BINDIR"

            BACKUP_DIR="/data/er_dashboard/backups"
            mkdir -p "$BACKUP_DIR"
            TIMESTAMP=$(date +%Y%m%d_%H%M%S)
            bashio::log.info "Creating safety backup before upgrade..."
            if tar -czf "$BACKUP_DIR/database_pre_upgrade_${TIMESTAMP}.tar.gz" -C /data/er_dashboard database 2>/dev/null; then
                bashio::log.info "Backup created successfully."
            else
                bashio::log.warning "Backup creation failed. Proceeding anyway..."
            fi

            rm -rf /data/er_dashboard/database_new
            mkdir -p /data/er_dashboard/database_new
            chown -R postgres:postgres /data/er_dashboard/database_new
            chmod 700 /data/er_dashboard/database_new

            echo "${POSTGRES_PASSWORD}" > /tmp/pgpass
            chown postgres:postgres /tmp/pgpass

            if ! su-exec postgres initdb -D /data/er_dashboard/database_new --username=postgres --auth=md5 --pwfile=/tmp/pgpass 2>&1; then
                bashio::log.error "initdb for upgrade failed!"
                rm -f /tmp/pgpass
                exit 1
            fi
            rm -f /tmp/pgpass

            bashio::log.info "Running pg_upgrade..."
            if (cd /tmp && su-exec postgres pg_upgrade \
                -b "$OLD_BINDIR" \
                -B "$NEW_BINDIR" \
                -d /data/er_dashboard/database \
                -D /data/er_dashboard/database_new); then

                bashio::log.info "pg_upgrade completed successfully!"
                rm -rf /data/er_dashboard/database
                mv /data/er_dashboard/database_new /data/er_dashboard/database
                bashio::log.info "✅ Database successfully upgraded to version $NEW_VERSION"
            else
                bashio::log.error "pg_upgrade failed! Check logs or output above."
                rm -rf /data/er_dashboard/database_new
                exit 1
            fi
        elif [ "$DB_VERSION" -gt "$NEW_VERSION" ]; then
            bashio::log.error "Database version ($DB_VERSION) is newer than PostgreSQL server ($NEW_VERSION)!"
            exit 1
        fi
    fi
fi
