# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/svc-paperless-web/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

export PAPERLESS_DATA_DIR=/site/data
export PAPERLESS_MEDIA_ROOT=/site/media
export PAPERLESS_CONSUMPTION_DIR=/site/consume
export PAPERLESS_REDIS="redis://localhost:6379"
export PAPERLESS_DBENGINE=mariadb
export PAPERLESS_DBHOST=$(bashio::services.mysql.host)
export PAPERLESS_DBPORT=$(bashio::services.mysql.port)
export PAPERLESS_DBNAME=$(bashio::services.mysql.database)
export PAPERLESS_DBUSER=$(bashio::services.mysql.username)
export PAPERLESS_DBPASS=$(bashio::services.mysql.password)
export PAPERLESS_URL="$(bashio::config 'url')"
export PAPERLESS_SECRET_KEY="$(bashio::config 'secret_key')"
export PAPERLESS_TIME_ZONE="$(bashio::config 'time_zone')"
export PAPERLESS_OCR_LANGUAGE="$(bashio::config 'ocr_language')"
export PAPERLESS_FILENAME_FORMAT="$(bashio::config 'filename_format')"

bashio::log.info "Starting Paperless Webserver..."
cd /usr/src/paperless
exec gunicorn \
    --bind 0.0.0.0:8000 \
    --workers 2 \
    paperless.wsgi:application
