# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/init-paperless/up
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

export PAPERLESS_DATA_DIR=/site/data
export PAPERLESS_MEDIA_ROOT=/site/media
export PAPERLESS_CONSUMPTION_DIR=/site/consume
export PAPERLESS_REDIS="redis://localhost:6379"

# Create directories if they don't exist
mkdir -p "$PAPERLESS_DATA_DIR" "$PAPERLESS_MEDIA_ROOT" "$PAPERLESS_CONSUMPTION_DIR"

bashio::log.info "Running database migrations..."
python3 /usr/src/paperless/manage.py migrate

# Create superuser if configured
if bashio::config.has_value 'admin_user' && bashio::config.has_value 'admin_password'; then
    export DJANGO_SUPERUSER_USERNAME=$(bashio::config 'admin_user')
    export DJANGO_SUPERUSER_PASSWORD=$(bashio::config 'admin_password')
    export DJANGO_SUPERUSER_EMAIL="${DJANGO_SUPERUSER_USERNAME}@example.com"

    bashio::log.info "Creating admin user: $DJANGO_SUPERUSER_USERNAME"
    if ! python3 /usr/src/paperless/manage.py createsuperuser --noinput; then
        bashio::log.warning "Failed to create superuser (it usually means it already exists)."
    fi
fi

bashio::log.info "Paperless-ngx initialization complete."
