# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/init-paperless/up
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "Waiting for MySQL service..."
if ! bashio::services.available 'mysql'; then
    bashio::log.error "MySQL service is not available!"
    bashio::exit.nok "MySQL service required"
fi

export PAPERLESS_DATA_DIR=/site/data
export PAPERLESS_MEDIA_ROOT=/site/media
export PAPERLESS_CONSUMPTION_DIR=/site/consume
export PAPERLESS_REDIS="redis://localhost:6379"
export PAPERLESS_DBENGINE=mariadb
export PAPERLESS_DBHOST=$(bashio::services.mysql.host)
export PAPERLESS_DBPORT=$(bashio::services.mysql.port)
export PAPERLESS_DBNAME=$(bashio::services.mysql.database)
export PAPERLESS_DBUSER=$(bashio::services.mysql.username)
export PAPERLESS_DBPASS=$(bashio::services.mysql.password)
export PAPERLESS_URL="$(bashio::config 'url')"
export PAPERLESS_SECRET_KEY="$(bashio::config 'secret_key')"
export PAPERLESS_TIME_ZONE="$(bashio::config 'time_zone')"
export PAPERLESS_OCR_LANGUAGE="$(bashio::config 'ocr_language')"
export PAPERLESS_FILENAME_FORMAT="$(bashio::config 'filename_format')"

bashio::log.info "MySQL connection: ${PAPERLESS_DBUSER}@${PAPERLESS_DBHOST}:${PAPERLESS_DBPORT}/${PAPERLESS_DBNAME}"

# Create directories if they don't exist
mkdir -p "$PAPERLESS_DATA_DIR" "$PAPERLESS_MEDIA_ROOT" "$PAPERLESS_CONSUMPTION_DIR"

# Database Reset
if bashio::config.true 'reset_database'; then
    if ! bashio::config.true 'reset_database_confirm'; then
        bashio::log.error "Database reset requires 'reset_database_confirm' to be true!"
        bashio::exit.nok "Database reset confirmation required"
    fi

    bashio::log.warning "=========================================="
    bashio::log.warning "   ⚠️  DATABASE RESET ENABLED  ⚠️"
    bashio::log.warning "=========================================="
    bashio::log.warning "ALL DATA WILL BE DELETED in 5 seconds!"
    sleep 5

    # Drop all tables (mariadb-cleaner approach or python loop)
    # Since we have mariadb-client (mariadb-dev), we can use mysql command if available or python
    # Using python to stay generic or just drop database if user has rights?
    # Usually addon user has one DB. dropping tables is safer.

    bashio::log.info "Dropping all tables..."
    # We can use a simple python script with django connection to flush
    python3 /usr/src/paperless/manage.py flush --no-input || bashio::log.warning "Failed to flush database"

    rm -f "$PAPERLESS_DATA_DIR/admin_created"

    # Reset options after successful reset (non-fatal if function doesn't exist)
    set +e
    bashio::addon.option 'reset_database' 2>/dev/null || true
    bashio::addon.option 'reset_database_confirm' 2>/dev/null || true
    set -e

    bashio::log.info "Database reset complete."
fi

bashio::log.info "Running database migrations..."
if ! python3 /usr/src/paperless/manage.py migrate; then
    bashio::log.error "Database migrations failed!"
    bashio::exit.nok "Failed to run database migrations"
fi

# Create superuser if configured
if ! bashio::fs.file_exists "$PAPERLESS_DATA_DIR/admin_created"; then
    bashio::log.info "First start detected - creating admin user..."

    ADMIN_USER=$(bashio::config 'admin_user')
    ADMIN_PASSWORD=$(bashio::config 'admin_password' '')
    if [ -z "${ADMIN_PASSWORD}" ]; then
        ADMIN_PASSWORD=$(openssl rand -base64 24)
    fi
    export DJANGO_SUPERUSER_USERNAME="$ADMIN_USER"
    export DJANGO_SUPERUSER_PASSWORD="$ADMIN_PASSWORD"
    export DJANGO_SUPERUSER_EMAIL="${ADMIN_USER}@example.com"

    if python3 /usr/src/paperless/manage.py createsuperuser --noinput; then
        touch "$PAPERLESS_DATA_DIR/admin_created"

        bashio::log.warning "=========================================="
        bashio::log.warning "   🔐 PAPERLESS ADMIN CREDENTIALS"
        bashio::log.warning "=========================================="
        bashio::log.warning "Username: $ADMIN_USER"
        bashio::log.warning "Password: $ADMIN_PASSWORD"
        bashio::log.warning "=========================================="
        bashio::log.warning "SAVE THIS PASSWORD NOW! It will NOT be shown again."
        bashio::log.warning "=========================================="
    else
        bashio::log.warning "Failed to create superuser."
    fi
fi

bashio::log.info "Paperless-ngx initialization complete."
