# vaultwarden/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/env bash
source /usr/lib/bashio/bashio.sh

# ==============================================================================
# Home Assistant Community Add-on: Vaultwarden
# Configures and runs Nginx
# ==============================================================================

if bashio::config.true 'ssl'; then
    certfile=$(bashio::config 'certfile')
    keyfile=$(bashio::config 'keyfile')
    cp /etc/nginx/servers/direct-ssl.disabled /etc/nginx/servers/direct.conf
    sed -i "s#%%certfile%%#${certfile}#g" /etc/nginx/servers/direct.conf
    sed -i "s#%%keyfile%%#${keyfile}#g" /etc/nginx/servers/direct.conf
else
    cp /etc/nginx/servers/direct.disabled /etc/nginx/servers/direct.conf
fi

bashio::log.info "Starting NGinx..."
exec nginx
