# vaultwarden/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Vaultwarden
# NGINX service
# ==============================================================================

bashio::log.info "Starting NGINX..."
exec nginx -g "daemon off;"
