# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/svc-backend/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Backend..."

# Application configuration
export PROJECT_NAME="ER-Startseite"
export BACKEND_CORS_ORIGINS='["http://localhost","http://127.0.0.1"]'
export POSTGRES_SERVER="localhost"
export POSTGRES_USER="postgres"
export POSTGRES_DB="er_startseite"

# Production settings for Uvicorn
# We use a single worker since we're using an embedded PostgreSQL instance
# which might not handle multiple connections from multiple processes well without a pooler.
if ! LOG_LEVEL=$(bashio::config 'log_level') || [ -z "$LOG_LEVEL" ]; then
    bashio::log.warning "Failed to fetch log_level configuration. Using default: info"
    LOG_LEVEL="info"
fi

# Map schema log levels to Uvicorn levels
# Schema supports: trace, debug, info, notice, warning, error, fatal
# Uvicorn supports: trace, debug, info, warning, error, critical
case "$LOG_LEVEL" in
    trace)   LOG_LEVEL="trace" ;;
    debug)   LOG_LEVEL="debug" ;;
    info)    LOG_LEVEL="info" ;;
    notice)  LOG_LEVEL="info" ;;   # Map notice → info
    warning) LOG_LEVEL="warning" ;;
    error)   LOG_LEVEL="error" ;;
    fatal)   LOG_LEVEL="critical" ;; # Map fatal → critical
    *)       LOG_LEVEL="info" ;;
esac

# Validate result
if [[ ! "$LOG_LEVEL" =~ ^(trace|debug|info|warning|error|critical)$ ]]; then
    bashio::log.error "Invalid log level '$LOG_LEVEL' for Uvicorn"
    bashio::exit.nok
fi

cd /app
exec python3 -m uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 13000 \
    --workers 1 \
    --log-level "$LOG_LEVEL" \
    --proxy-headers \
    --forwarded-allow-ips='*'
