# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/svc-backend/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Backend..."

# Application configuration
export PROJECT_NAME="ER-Startseite"
export BACKEND_CORS_ORIGINS='["http://localhost","http://127.0.0.1"]'
export POSTGRES_SERVER="localhost"
export POSTGRES_USER="postgres"
export POSTGRES_DB="er_startseite"

# Production settings for Uvicorn
# We use a single worker since we're using an embedded PostgreSQL instance
# which might not handle multiple connections from multiple processes well without a pooler.
LOG_LEVEL="info"
if bashio::config.has_value 'log_level'; then
    LOG_LEVEL=$(bashio::config 'log_level')
fi

cd /app
exec python3 -m uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 13000 \
    --workers 1 \
    --log-level "$LOG_LEVEL" \
    --proxy-headers \
    --forwarded-allow-ips='*'
