# netboot-xyz/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

# Map Bashio log level to Nginx log level
log_level=$(bashio::config 'log_level')
case "${log_level}" in
    trace) nginx_log_level="debug" ;;
    debug) nginx_log_level="debug" ;;
    info)  nginx_log_level="info" ;;
    notice) nginx_log_level="notice" ;;
    warning) nginx_log_level="warn" ;;
    error) nginx_log_level="error" ;;
    fatal) nginx_log_level="crit" ;;
    *)     nginx_log_level="error" ;;
esac

sed -i "s|error_log .*|error_log /var/log/nginx/error.log $nginx_log_level;|" /config/nginx/nginx.conf 2>/dev/null || true
# Also update main nginx.conf if exists
sed -i "s|error_log .*|error_log /var/log/nginx/error.log $nginx_log_level;|" /etc/nginx/nginx.conf 2>/dev/null || true

exec /usr/sbin/nginx -c /config/nginx/nginx.conf