# ShieldDNS/rootfs/etc/services.d/coredns/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Port configured. Starting CoreDNS..."
exec /usr/bin/coredns -conf /etc/Corefile
