# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Nginx..."

# Map Bashio log_level to Nginx log_level
log_level=$(bashio::config 'log_level')
nginx_log_level="warn"
case "${log_level}" in
    trace|debug) nginx_log_level="debug" ;;
    info)        nginx_log_level="info" ;;
    notice)      nginx_log_level="notice" ;;
    warning)     nginx_log_level="warn" ;;
    error)       nginx_log_level="error" ;;
    fatal)       nginx_log_level="crit" ;;
    *)           nginx_log_level="warn" ;;
esac

# Generate nginx config from template
if ! bashio::var.json \
    ssl "$(bashio::config 'ssl')" \
    certfile "$(bashio::config 'certfile')" \
    keyfile "$(bashio::config 'keyfile')" \
    log_level "${nginx_log_level}" \
    | tempio \
        -template /etc/nginx/templates/direct.gtpl \
        -out /etc/nginx/http.d/direct.conf; then
    bashio::log.error "Failed to generate nginx configuration from template!"
    exit 1
fi

# Generate Ingress config
if bashio::config.true 'ingress'; then
    bashio::log.info "Generating Ingress configuration..."
    if ! bashio::var.json \
        log_level "${nginx_log_level}" \
        | tempio \
            -template /etc/nginx/templates/ingress.gtpl \
            -out /etc/nginx/http.d/ingress.conf; then
        bashio::log.error "Failed to generate ingress configuration!"
    fi
fi

# Validate SSL files if SSL is enabled
if bashio::config.true 'ssl'; then
    certfile="/ssl/$(bashio::config 'certfile')"
    keyfile="/ssl/$(bashio::config 'keyfile')"

    if [[ ! -f "$certfile" ]]; then
        bashio::log.error "SSL certificate not found: $certfile"
        exit 1
    fi
    if [[ ! -f "$keyfile" ]]; then
        bashio::log.error "SSL key not found: $keyfile"
        exit 1
    fi
    bashio::log.info "SSL certificates validated."
fi

exec nginx -g "daemon off;"
