# er-dashboard/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Nginx..."

bashio::var.json \
    ssl "$(bashio::config 'ssl')" \
    certfile "$(bashio::config 'certfile')" \
    keyfile "$(bashio::config 'keyfile')" \
    | tempio \
        -template /etc/nginx/templates/direct.gtpl \
        -out /etc/nginx/http.d/direct.conf

exec nginx -g "daemon off;"
