# antigravity-server/rootfs/etc/s6-overlay/s6-rc.d/vnc-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: Antigravity-Server
# Runs the TigerVNC server
# ==============================================================================

bashio::log.info "Starting TigerVNC server on display :1..."

# Start VNC server
exec vncserver :1 \
    -localhost no \
    -geometry 1920x1080 \
    -depth 24 \
    -fg \
    -SecurityTypes VncAuth \
    -rfbauth /root/.config/tigervnc/passwd
