# addons/n8n/rootfs/etc/s6-overlay/s6-rc.d/n8n-init/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: N8n
# Configures N8n before the main service starts
# ==============================================================================
bashio::log.info "Configuring N8n..."

# Ensure persistence
if ! bashio::fs.directory_exists "/data/n8n"; then
    bashio::log.info "Creating persistent storage directory..."
    mkdir -p /data/n8n
fi

# Note: Ownership is not explicitly set here because the add-on runs as root
# in the base image and /data is managed by the Supervisor. If running as a
# non-root user in the future, uncomment and adjust the following line:
# chown -R n8n:n8n /data/n8n

bashio::log.info "Configuration complete."
