# ShieldFile/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC1091

bashio::log.info "🚀 Starting Nginx..."

# Configure Nginx log level
if ! LOG_LEVEL=$(bashio::config 'log_level') || [ -z "$LOG_LEVEL" ]; then
    bashio::log.warning "Failed to fetch log_level configuration. Using default: info"
    LOG_LEVEL="info"
fi
# Map Bashio levels to Nginx levels (debug, info, notice, warn, error, crit, alert, emerg)
case "${LOG_LEVEL}" in
        trace|debug) NGINX_LOG_LEVEL="debug" ;;
        info)        NGINX_LOG_LEVEL="info" ;;
        notice)      NGINX_LOG_LEVEL="notice" ;;
        warning)     NGINX_LOG_LEVEL="warn" ;;
        error)       NGINX_LOG_LEVEL="error" ;;
        fatal)       NGINX_LOG_LEVEL="crit" ;;
        *)           NGINX_LOG_LEVEL="warn" ;;
    esac
    sed -i -E "s/error_log .*/error_log \/dev\/stderr ${NGINX_LOG_LEVEL};/g" /etc/nginx/nginx.conf


exec nginx -g 'daemon off;'
