# switchcraft/rootfs/etc/s6-overlay/s6-rc.d/switchcraft/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: SwitchCraft
# Service definition
# ==============================================================================
set -x
LOG_LEVEL=$(bashio::config 'log_level')
FLET_SECRET=$(bashio::config 'flet_secret' 'null')

bashio::log.info "Starting secret handling..."
# Handle FLET_SECRET
if [ "${FLET_SECRET}" = "null" ] || [ -z "${FLET_SECRET}" ]; then
    if [ -f "/data/flet_secret" ]; then
        FLET_SECRET=$(cat /data/flet_secret)
        bashio::log.info "Using persisted FLET_SECRET"
    else
        bashio::log.info "Generating new FLET_SECRET"
        FLET_SECRET=$(tr -dc 'a-zA-Z0-9' </dev/urandom | head -c 32)
        ( umask 077; echo "${FLET_SECRET}" > /data/flet_secret )
    fi
fi
export FLET_SECRET

cd /app || bashio::exit.nok "Could not change to /app directory"
bashio::log.info "Changed to /app"
export PYTHONPATH="/app/src:${PYTHONPATH}"
export HOME=/data
bashio::log.info "Exported env vars"

# SSO Environment Variables
if bashio::config.has_value 'entra_client_id'; then
    export SC_ENTRA_CLIENT_ID=$(bashio::config 'entra_client_id')
    export SC_ENTRA_CLIENT_SECRET=$(bashio::config 'entra_client_secret')
    export SC_ENTRA_TENANT_ID=$(bashio::config 'entra_tenant_id')
fi

if bashio::config.has_value 'github_client_id'; then
    export SC_GITHUB_CLIENT_ID=$(bashio::config 'github_client_id')
    export SC_GITHUB_CLIENT_SECRET=$(bashio::config 'github_client_secret')
fi

if bashio::config.has_value 'base_url'; then
    export SC_BASE_URL=$(bashio::config 'base_url')
fi

bashio::log.info "Starting SwitchCraft (Auth Proxy Mode)..."
exec python3 -m uvicorn switchcraft.server.app:app \
    --host 0.0.0.0 \
    --port 8080 \
    --proxy-headers \
    --forwarded-allow-ips='*' \
    --log-level "${LOG_LEVEL,,}"