# switchcraft/rootfs/etc/s6-overlay/s6-rc.d/switchcraft/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App: SwitchCraft
# Service definition
# ==============================================================================
set -e

# Trace execution for debugging
# Trace execution for debugging (Disabled)
# set -x

bashio::log.info "Starting SwitchCraft App Service..."

# Get Config
LOG_LEVEL=$(bashio::config 'log_level')
bashio::log.info "Log Level: ${LOG_LEVEL}"

# handle flet_secret
if bashio::config.has_value 'flet_secret'; then
    bashio::log.info "Using configured FLET_SECRET"
    FLET_SECRET=$(bashio::config 'flet_secret')
else
    bashio::log.info "FLET_SECRET not provided in config"

    if [ -f "/data/flet_secret" ]; then
        bashio::log.info "Using persisted FLET_SECRET from /data/flet_secret"
        FLET_SECRET=$(cat /data/flet_secret)
    else
        bashio::log.info "Generating new FLET_SECRET..."
        # Safer generation using openssl which is usually faster and cleaner than urandom pipe
        if command -v openssl >/dev/null 2>&1; then
             FLET_SECRET=$(openssl rand -base64 32 | tr -dc 'a-zA-Z0-9' | head -c 32)
        else
             FLET_SECRET=$(tr -dc 'a-zA-Z0-9' < /dev/urandom | head -c 32)
        fi

        bashio::log.info "Saving FLET_SECRET..."
        echo "${FLET_SECRET}" > /data/flet_secret
        chmod 600 /data/flet_secret
    fi
fi

export FLET_SECRET

bashio::log.info "Changing directory to /app..."
cd /app || bashio::exit.nok "Could not change to /app directory"

# Safely handle PYTHONPATH
if [ -z "${PYTHONPATH:-}" ]; then
    export PYTHONPATH="/app/src"
else
    export PYTHONPATH="/app/src:${PYTHONPATH}"
fi
export HOME=/data

# Export optional env vars
if bashio::config.has_value 'entra_client_id'; then
    export SC_ENTRA_CLIENT_ID=$(bashio::config 'entra_client_id')
    export SC_ENTRA_CLIENT_SECRET=$(bashio::config 'entra_client_secret')
    export SC_ENTRA_TENANT_ID=$(bashio::config 'entra_tenant_id')
fi

if bashio::config.has_value 'github_client_id'; then
    export SC_GITHUB_CLIENT_ID=$(bashio::config 'github_client_id')
    export SC_GITHUB_CLIENT_SECRET=$(bashio::config 'github_client_secret')
fi

if bashio::config.has_value 'base_url'; then
    export SC_BASE_URL=$(bashio::config 'base_url')
fi

bashio::log.info "Starting uvicorn server on port 8080..."

# Execute uvicorn
exec python3 -m uvicorn switchcraft.server.app:app \
    --host 0.0.0.0 \
    --port 8080 \
    --proxy-headers \
    --forwarded-allow-ips='*' \
    --log-level "${LOG_LEVEL,,}"
