# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/svc-postgresql/run
#!/command/with-contenv bashio

bashio::log.info "Starting PostgreSQL..."

exec su-exec postgres postgres -D /data/database
