# tt-rss/rootfs/etc/s6-overlay/s6-rc.d/svc-postgresql/run
#!/usr/bin/env bashio

bashio::log.info "Starting PostgreSQL..."

exec s6-setuidgid postgres postgres -D /data/database
