# whatsapp/rootfs/etc/s6-overlay/s6-rc.d/whatsapp/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting WhatsApp Addon (Baileys)..."

export PORT=8066
# Map Bashio log level to Pino log level (used by Baileys/WhatsApp)
if ! LOG_LEVEL=$(bashio::config 'log_level') || [ -z "$LOG_LEVEL" ]; then
    bashio::log.warning "Failed to fetch log_level configuration. Using default: info"
    LOG_LEVEL="info"
fi
case "${LOG_LEVEL}" in
    trace)        LOG_LEVEL="silly" ;;
    debug)        LOG_LEVEL="debug" ;;
    info|notice)  LOG_LEVEL="info"  ;;
    warning|warn) LOG_LEVEL="warn"  ;;
    error|fatal)  LOG_LEVEL="error" ;;
    *)            LOG_LEVEL="info"  ;;
esac
export LOG_LEVEL
if bashio::config.true 'reset_session'; then
    bashio::log.warning "⚠️  RESET_SESSION toggle is ENABLED."
    bashio::log.warning "Clearing authentication data in /data/auth_info_baileys..."
    rm -rf /data/auth_info_baileys
    export RESET_SESSION="true"
else
    export RESET_SESSION="false"
fi

# Send Message Timeout (Default 25s)
if ! SEND_MESSAGE_TIMEOUT=$(bashio::config 'send_message_timeout') || [ -z "$SEND_MESSAGE_TIMEOUT" ]; then
    SEND_MESSAGE_TIMEOUT=25000
fi
export SEND_MESSAGE_TIMEOUT

# Keep Alive Interval (Default 30s)
if ! KEEP_ALIVE_INTERVAL=$(bashio::config 'keep_alive_interval') || [ -z "$KEEP_ALIVE_INTERVAL" ]; then
    KEEP_ALIVE_INTERVAL=30000
fi
export KEEP_ALIVE_INTERVAL

# Mask Sensitive Data (Default false)
if bashio::config.true 'mask_sensitive_data'; then
    export MASK_SENSITIVE_DATA="true"
else
    export MASK_SENSITIVE_DATA="false"
fi

# Webhook Configuration
if bashio::config.true 'webhook_enabled'; then
    export WEBHOOK_ENABLED="true"
else
    export WEBHOOK_ENABLED="false"
fi

if ! WEBHOOK_URL=$(bashio::config 'webhook_url') || [ "$WEBHOOK_URL" = "null" ] || [ -z "$WEBHOOK_URL" ]; then
    WEBHOOK_URL=""
fi
export WEBHOOK_URL

if ! WEBHOOK_TOKEN=$(bashio::config 'webhook_token') || [ "$WEBHOOK_TOKEN" = "null" ] || [ -z "$WEBHOOK_TOKEN" ]; then
    WEBHOOK_TOKEN=""
fi
export WEBHOOK_TOKEN

# UI Authentication
if bashio::config.true 'ui_auth_enabled'; then
    export UI_AUTH_ENABLED="true"
else
    export UI_AUTH_ENABLED="false"
fi

if ! UI_AUTH_PASSWORD=$(bashio::config 'ui_auth_password') || [ "$UI_AUTH_PASSWORD" = "null" ] || [ -z "$UI_AUTH_PASSWORD" ]; then
    UI_AUTH_PASSWORD=""
fi
export UI_AUTH_PASSWORD

# Mark Online (Default false)
if bashio::config.true 'mark_online'; then
    export MARK_ONLINE="true"
else
    export MARK_ONLINE="false"
fi

# Media Folder Configuration
if ! MEDIA_FOLDER=$(bashio::config 'media_folder') || [ "$MEDIA_FOLDER" = "null" ] || [ -z "$MEDIA_FOLDER" ]; then
    MEDIA_FOLDER=""
fi
export MEDIA_FOLDER

# Detect Home Assistant config directory
HA_CONFIG_ROOT=""
for path in "/config" "/homeassistant"; do
    if [ -d "$path" ]; then
        HA_CONFIG_ROOT="$path"
        break
    fi
done

if [ -z "$HA_CONFIG_ROOT" ]; then
    bashio::log.error "Could not find Home Assistant configuration directory in /config or /homeassistant"
    # Fallback to /config and hope for the best
    HA_CONFIG_ROOT="/config"
fi

INTEGRATION_DIR="${HA_CONFIG_ROOT}/custom_components/whatsapp"

# Function to compare semantic versions
version_gt() {
    test "$(printf '%s\n' "$@" | sort -V | head -n 1)" != "$1"
}

# Function to install/update integration
install_integration() {
    local TAG_NAME=$1
    local IS_UPDATE=$2

    bashio::log.info "Installing integration version: ${TAG_NAME:-Default Branch}..."

    if [ -n "$TAG_NAME" ]; then
        CLONE_ARGS=("--branch" "$TAG_NAME")
    else
        CLONE_ARGS=()
    fi

    mkdir -p "/tmp/ha-whatsapp_install"
    rm -rf "/tmp/ha-whatsapp_install" # Ensure clean state

    if git clone --depth 1 "${CLONE_ARGS[@]}" https://github.com/FaserF/ha-whatsapp.git "/tmp/ha-whatsapp_install" > /dev/null 2>&1; then
        if [ -d "/tmp/ha-whatsapp_install/custom_components/whatsapp" ]; then
            # Check if we are updating, if so remove old
            if [ -d "$INTEGRATION_DIR" ]; then
                 bashio::log.info "Removing old integration version..."
                 rm -rf "$INTEGRATION_DIR"
            fi

            mkdir -p "${HA_CONFIG_ROOT}/custom_components"
            # Use cp -rf for forceful copy
            cp -rf "/tmp/ha-whatsapp_install/custom_components/whatsapp" "${HA_CONFIG_ROOT}/custom_components/"

            # Double check installation
            if [ -d "$INTEGRATION_DIR" ]; then
                bashio::log.info "✅ Integration successfully installed/updated to $INTEGRATION_DIR"

                if [ "$IS_UPDATE" = "true" ]; then
                    MSG="The WhatsApp integration has been updated to $TAG_NAME. Please restart Home Assistant."
                    TITLE="WhatsApp Integration Updated"
                else
                    MSG="The WhatsApp integration has been installed. Please restart Home Assistant."
                    TITLE="WhatsApp Integration Installed"
                fi

                bashio::log.warning "RESTART HOME ASSISTANT to load the changes!"
                bashio::log.info "Publishing notification to Home Assistant..."

                curl -s -X POST \
                    -H "Authorization: Bearer $SUPERVISOR_TOKEN" \
                    -H "Content-Type: application/json" \
                    -d "{\"title\": \"$TITLE\", \"message\": \"$MSG\", \"notification_id\": \"whatsapp_restart_required\"}" \
                    http://supervisor/core/api/services/persistent_notification/create > /dev/null
            else
                bashio::log.error "❌ Copy failed: $INTEGRATION_DIR still does not exist."
            fi
        else
            bashio::log.error "❌ Could not find 'custom_components/whatsapp' in the repository."
        fi
    else
        bashio::log.error "❌ Failed to clone repository."
    fi
    rm -rf "/tmp/ha-whatsapp_install"
}

# Determine Channel (Edge vs Stable)
# We check slug, name, and version to be sure
# Fetch Addon Info from Supervisor API directly check
# (bashio::addon.* commands were failing, likely due to older bashio version in base or context)
ADDON_INFO=$(curl -s -H "Authorization: Bearer $SUPERVISOR_TOKEN" http://supervisor/addons/self/info)
SLUG=$(echo "$ADDON_INFO" | jq -r '.data.slug // empty')
NAME=$(echo "$ADDON_INFO" | jq -r '.data.name // empty')
VERSION=$(echo "$ADDON_INFO" | jq -r '.data.version // empty')

bashio::log.info "Channel Detection - Slug: ${SLUG}, Name: ${NAME}, Version: ${VERSION}"

# Check if any indicator points to Edge/Dev channel
# 1. Slug contains edge
# 2. Name contains Edge (case-insensitive)
# 3. Version contains 'dev', 'git', or a commit-like hash (e.g. 7+ characters)
if [[ "${SLUG}" == *"edge"* ]] || \
   [[ "${NAME,,}" == *"edge"* ]] || \
   [[ "${VERSION}" == *"dev"* ]] || \
   [[ "${VERSION}" == *"git"* ]] || \
   [[ "${VERSION}" =~ [0-9a-f]{7,40} ]]; then
    CHANNEL="edge"
    bashio::log.info "🟢 Edge/Dev channel detected. Will prefer Pre-releases."
else
    CHANNEL="stable"
    bashio::log.info "🔵 Stable channel detected."
fi

# Fetch Releases once
bashio::log.info "Fetching release information..."
RELEASES_JSON=$(curl -s -A "HomeAssistant-Addon" "https://api.github.com/repos/FaserF/ha-whatsapp/releases")

# Valid JSON check
if ! echo "$RELEASES_JSON" | jq -e 'if type=="array" then true else false end' > /dev/null 2>&1; then
    bashio::log.info "Note: Could not fetch releases (GitHub API may be rate-limited). Skipping update checks."
    RELEASES_JSON="[]"
fi

# Determine Latest Versions based on channel
TARGET_TAG=""
if [ "$RELEASES_JSON" != "[]" ]; then
    if [ "$CHANNEL" == "edge" ]; then
        # Edge: Latest release (whether stable or beta)
        TARGET_TAG=$(echo "$RELEASES_JSON" | jq -r '.[0].tag_name // empty')
    else
        # Stable: Latest non-prerelease
        TARGET_TAG=$(echo "$RELEASES_JSON" | jq -r 'map(select(.prerelease == false)) | .[0].tag_name // empty')
    fi

    # Fallback: If stable search failed but we have releases, ensure we dont leave TARGET_TAG empty if possible
    if [ -z "$TARGET_TAG" ] && [ "$CHANNEL" == "stable" ]; then
        if [ "$(echo "$RELEASES_JSON" | jq 'length')" -gt 0 ]; then
            fallback=$(echo "$RELEASES_JSON" | jq -r '.[0].tag_name')
            bashio::log.info "No stable release found, falling back to latest tag: $fallback"
            TARGET_TAG="$fallback"
        fi
    fi
fi


# Manage Integration
if [ ! -d "$INTEGRATION_DIR" ]; then
    bashio::log.warning "Integration not found in $INTEGRATION_DIR"

    if [ -n "$TARGET_TAG" ]; then
        bashio::log.info "Target version (Initial Install): $TARGET_TAG"
        install_integration "$TARGET_TAG" "false"
    else
        bashio::log.info "No releases found. Installing from Main branch..."
        install_integration "" "false"
    fi
else
    bashio::log.info "✅ Integration found at $INTEGRATION_DIR"
    bashio::log.info "Checking for updates..."

    if [ -f "$INTEGRATION_DIR/manifest.json" ]; then
        CURRENT_VERSION=$(jq -r '.version' "$INTEGRATION_DIR/manifest.json")
        bashio::log.info "Current Version: $CURRENT_VERSION"

        UPDATE_NEEDED="false"

        # Comparison Logic
        if [ -n "$TARGET_TAG" ]; then
             # We should compare full strings, but often manifests are 0.0.1 and tags are v0.0.1
             # Let's normalize
             curr=$(echo "$CURRENT_VERSION" | sed 's/^v//')
             targ=$(echo "$TARGET_TAG" | sed 's/^v//')

             if [ "$curr" != "$targ" ]; then
                 # If strings differ, check semantic
                 if version_gt "$targ" "$curr"; then
                     bashio::log.info "Update Available: $targ > $curr"
                     UPDATE_NEEDED="true"
                 else
                     bashio::log.info "Current version ($curr) is equal or newer than target ($targ)."
                 fi
             else
                 bashio::log.info "Versions match ($curr)."
             fi
        elif [ "$RELEASES_JSON" != "[]" ]; then
             bashio::log.warning "Could not determine target tag from GitHub releases."
        fi

        if [ "$UPDATE_NEEDED" == "true" ]; then
             bashio::log.info "⬆️  Auto-updating to $TARGET_TAG..."
             install_integration "$TARGET_TAG" "true"
        else
             bashio::log.info "✨ Integration is up to date."
        fi
    else
        bashio::log.warning "Manifest not found. Cannot check version. Reinstalling..."
        install_integration "$TARGET_TAG" "true"
    fi
fi

cd /opt/whatsapp
exec node index.js
