# whatsapp/rootfs/etc/s6-overlay/s6-rc.d/whatsapp/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting WhatsApp Addon (Baileys)..."

cd /opt/whatsapp
exec npm start
