# whatsapp/rootfs/etc/s6-overlay/s6-rc.d/whatsapp/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting WhatsApp Addon (Baileys)..."

export PORT=$(bashio::addon.port 8066)
export MDNS_NAME=$(bashio::config 'mdns_name')

cd /opt/whatsapp
exec node index.js
