# komodo/rootfs/etc/s6-overlay/s6-rc.d/mongo/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "Starting MongoDB..."

# Ensure data directory exists
if [ ! -d /data/db ]; then
    mkdir -p /data/db
    chown mongodb:mongodb /data/db
fi

exec sudo -u mongodb mongod --dbpath /data/db --bind_ip 127.0.0.1
