# vaultwarden/rootfs/etc/s6-overlay/s6-rc.d/base-addon-banner/up
#!/usr/bin/env bash
# Disable strict mode to prevent crash in base image logic if repeated
set +u

# Minimal banner to satisfy the service requirement without crashing
echo " [Info] Starting Vaultwarden Add-on..."

# Exit successfully
exit 0
